# deploy/systemd/nginx.service.d/jts-recovery.conf
# JTS management front door hardening.
#
# nginx is package-owned, but JTS depends on it for the browser portal,
# static assets, and every nginx-fronted wizard. A transient OOM kill must
# not leave port 80 dark until an operator SSHes in.
[Unit]
StartLimitIntervalSec=600
StartLimitBurst=20

[Service]
Restart=always
RestartSec=5

# Protect the recovery surface, but keep it below jasper-control/voice/audio
# in the OOM ladder so the kernel still prefers killing web before the core
# speaker path.
OOMScoreAdjust=-450
