# deploy/debian-stack/systemd/go-librespot.service
[Unit]
Description=go-librespot Spotify Connect daemon (JTS)
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/go-librespot --config_dir /etc/go-librespot
User=pi
Group=audio
Restart=on-failure
RestartSec=2
MemoryHigh=150M
MemoryMax=200M
Nice=-5

[Install]
WantedBy=multi-user.target
