# deploy/systemd/bt-agent.service
[Unit]
Description=Jasper Bluetooth no-code pairing agent
Requires=bluetooth.service
After=bluetooth.service

[Service]
Type=simple
# JTS is a headless speaker. Its single default BlueZ agent advertises
# NoInputNoOutput so pairing is no-code "Just Works"; PIN/passkey/numeric
# comparison flows are rejected as unsupported.
ExecCondition=/opt/jasper/.venv/bin/jasper-local-source-allowed
ExecStart=/opt/jasper/.venv/bin/jasper-bluetooth-agent
TimeoutStopSec=10s
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
