# deploy/systemd/bt-agent.service
[Unit]
Description=Jasper Bluetooth no-code pairing agent
Requires=bluetooth.service
After=bluetooth.service
# Safe to restart after transient Bluetooth/OOM churn; finite so a broken
# BlueZ stack does not spin forever.
StartLimitIntervalSec=600
StartLimitBurst=20

[Service]
Type=simple
TimeoutStartSec=2s
# JTS is a headless speaker. Its single default BlueZ agent advertises
# NoInputNoOutput so pairing is no-code "Just Works"; PIN/passkey/numeric
# comparison flows are rejected as unsupported.
UnsetEnvironment=LD_PRELOAD LD_LIBRARY_PATH LD_AUDIT GLIBC_TUNABLES
ExecCondition=+/usr/bin/env -i PATH=/opt/jasper/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin /opt/jasper/.venv/bin/jasper-local-source-allowed --source bluetooth
ExecStart=/opt/jasper/.venv/bin/jasper-bluetooth-agent
TimeoutStopSec=10s
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
