# deploy/systemd/jasper-fanin-combo-health.timer
# Periodic tick for the USB-combo runtime-fallback watcher
# (jasper-fanin-combo-health.service). See docs/HANDOFF-usbsink.md.
[Unit]
Description=Periodically check whether the USB combo must be withdrawn and disarmed
Documentation=file:/opt/jasper/docs/HANDOFF-usbsink.md

# ~3-minute cadence, not seconds (mirrors jasper-wifi-recover). fan-in already
# self-heals a broken direct-capture handle within ~1-2 s via its own reopen; this
# watcher's unique job is the RARE case where that self-heal is not restoring
# durable capture, and the disarm requires >= 2 consecutive broken ticks (~6 min)
# so a single transient never trips the fallback — a few-minutes detection window
# is exactly right. A healthy tick produces NO script output (only a real
# transition logs), but each oneshot activation still costs ~2 systemd journal
# lines, the other reason the cadence is minutes not 90 s. OnBootSec is past the
# --auto boot pass (jasper-fanin-coupling-auto.service) so the first tick observes
# the freshly-armed combo rather than racing its arm.
[Timer]
OnBootSec=5min
OnUnitActiveSec=3min
AccuracySec=30s
RandomizedDelaySec=20s
Unit=jasper-fanin-combo-health.service
Persistent=false

[Install]
WantedBy=timers.target
