# deploy/systemd/jasper-aec-bridge.service
[Unit]
Description=Software AEC bridge — WebRTC AEC3 cancels speaker echo from XVF raw mic
# Captures the far-end reference from outputd's UDP final-speaker
# monitor and the near-end mic from hw:CARD=Array,DEV=0 channel 1 on
# the 6-ch firmware. In opt-in chip-AEC mode, the same outputd
# reference also feeds the XVF USB-IN chip reference path and the
# bridge forwards a chip ASR beam to
# 127.0.0.1:JASPER_AEC_UDP_PORT (default 9876) instead of running
# WebRTC AEC3.
#
# Historical note: this used to write to a second snd-aloop card
# (LoopbackAEC) that jasper-voice read via PortAudio. Replaced
# with UDP in PR 2 of the resilience-ladder series — snd-aloop's
# kernel-side loopback_cable wedged on consumer SIGKILL and
# required a reboot to clear. UDP has no such state.
#
# Engine: WebRTC AEC3 via the jasper_aec3 pybind11 binding around
# Trixie's libwebrtc-audio-processing-1 v1.3-3. Tuning knobs
# (JASPER_AEC_REF_GAIN_DB, JASPER_AEC_MIC_GAIN_DB, JASPER_AEC_AGC2)
# come from /etc/jasper/jasper.env — see docs/HANDOFF-aec.md
# "Tuning findings" for the chosen production values.
#
# REPLACES the previous CamillaDSP-based bridge. The default production
# path is WebRTC AEC3; the user-authorized chip-AEC carve-out is narrow
# and opt-in, with outputd owning the XVF USB-IN reference producer.
# See jasper/cli/aec_bridge.py and docs/CHIP-AEC-EXPERIMENT.md.
After=jasper-fanin.service jasper-camilla.service jasper-outputd.service jasper-aec-init.service sound.target
Requires=jasper-camilla.service
# `Wants=` so init runs whenever bridge starts (init is a one-shot
# that applies XVF tuning; bridge depends on chip being initialised).
# Soft dep — bridge still tries to start even if init failed; AEC
# quality degrades but the daemon is operational.
Wants=jasper-fanin.service jasper-outputd.service jasper-aec-init.service
# `PartOf=` propagates camilla's stop / restart to this unit. Without
# it, a camilla restart (e.g. during deploy) deactivates the bridge
# via `Requires=` cascade — but `Restart=on-failure` doesn't trigger
# on a clean exit (code 0 from cascade-stop), so the bridge stays
# dead silently. Mic ends up listening to a (now-unfed) loopback,
# wake events stop firing, speaker appears alive but is deaf. Hit
# in production on 2026-05-10.
PartOf=jasper-camilla.service

# Tier 4.5 — escalate to system reboot if the per-service watchdog
# spiral exceeds 4 restarts in 5 minutes. See
# docs/HANDOFF-tier5-watchdog-liveness.md "Option B (T5.1)" for
# rationale. `reboot` (not `reboot-force`): clean unmount needed
# on a 1 GB Pi to flush dirty zram pages.
StartLimitIntervalSec=300
StartLimitBurst=4
StartLimitAction=reboot

# Stage 2 audio-protection: assign to jts-mic.slice with
# MemorySwapMax=0. Today's stress test (2026-05-24) caught
# aec-bridge accumulating 42 MB of VmSwap during sustained
# memory pressure — exactly the failure mode this slice prevents.
# See docs/HANDOFF-resilience.md "Stage 2 — audio-protection
# subset (shipped)".

[Service]
# Tier 1+2 of the resilience ladder. `Type=notify` makes the daemon
# emit `READY=1` (via sdnotify, see jasper/watchdog.py) before
# systemd considers it active. `WatchdogSec=30s` requires the daemon
# to pat the watchdog every < 30 s; the daemon's heartbeat thread
# only pats if its work loop has made progress recently (sentinel
# pattern). If the loop wedges (PortAudio blocked in a syscall
# holding the GIL — the 2026-05-11 incident), `WATCHDOG=1` stops
# firing, this timer expires, and `Restart=on-watchdog` revives us
# with a fresh process. Critically, this happens BEFORE the SIGTERM
# timeout / SIGKILL that would otherwise corrupt snd-aloop's
# kernel-side loopback_cable state and require a reboot.
#
# `TimeoutStopSec=5s` is the matching shorter SIGTERM grace period
# (down from systemd's default 90 s). If the daemon can't shut down
# in 5 s, it's wedged and we want to get to SIGKILL+restart fast
# rather than accumulating wedge time.
Slice=jts-mic.slice
Type=notify
WatchdogSec=30s
TimeoutStopSec=5s
# Pull JASPER_AEC_REF_GAIN_DB, JASPER_AEC_MIC_GAIN_DB, JASPER_AEC_AGC2
# from the same env file jasper-voice uses. Leading `-` makes systemd
# skip the file silently if it's absent (e.g. fresh install before
# install.sh has run); the bridge defaults are conservative
# (ref_gain=0 dB, mic_gain=0 dB, agc2=off) so a missing env file is
# safe — but the production tuning values won't be applied.
EnvironmentFile=-/etc/jasper/jasper.env
# Corpus-recorder experiment flags live under /var/lib/jasper so the
# sandboxed web UI can deliberately enable optional bridge outputs
# (reference / cheap USB / DTLN) without writing to /etc.
EnvironmentFile=-/var/lib/jasper/wake_corpus_bridge.env
# Skip cleanly when the XVF Array is absent, on 2-channel firmware, or
# the operator has disabled AEC. The reconciler owns env/service state;
# this condition keeps manual starts from re-entering a crash loop.
ExecCondition=/usr/local/sbin/jasper-aec-reconcile --check-aec-ready
ExecStart=/opt/jasper/.venv/bin/jasper-aec-bridge
# `Restart=on-failure` covers non-zero exit + signals + watchdog
# timeouts (per the table in systemd.service(5)). The previous
# `on-watchdog` setting only restarted on watchdog timeout, so the
# explicit `BridgeStalled` exit (code 1, raised by _aec_loop when
# the mic queue stays empty for JASPER_AEC_STALL_RESTART_SEC) was
# NOT being restarted — the service was just parked as failed, and
# the wake-word path was offline until manual intervention. Hit in
# production 2026-05-11.
Restart=on-failure
RestartSec=2
Nice=-10
IOSchedulingClass=realtime
IOSchedulingPriority=4
LimitRTPRIO=99
LimitMEMLOCK=infinity

# OOM protection (Stage 1 of memory-resilience plan, see
# docs/HANDOFF-resilience.md). Real-time mic processing: killing
# this stops AEC and the wake loop runs raw, but the speaker can
# still play music — so we protect heavily but not as hard as
# camilla (-900). -700 keeps us well above any non-jasper process
# the kernel might consider.
OOMScoreAdjust=-700

[Install]
WantedBy=multi-user.target
