# deploy/systemd/jasper-fanin-combo-health.service
# RUNTIME-FALLBACK watcher for the USB low-latency combo (defect 2026-07-10).
# One oneshot tick of the coupling reconciler's --health verb: read fan-in's USB
# direct-capture health and, if it has been BROKEN for >= 2 consecutive ticks,
# disarm the combo back to the aloop bridge (the reconciler is the single writer,
# so it disarms exactly the way it arms) and write the flap-proof fallback marker.
# Healthy ticks are journal-quiet — only real transitions log. Fired by
# jasper-fanin-combo-health.timer on a ~3-min cadence (no resident daemon; mirrors
# jasper-wifi-recover). The combo is re-armed only by the --auto pass
# (jasper-fanin-coupling-auto.service) on boot/deploy/toggle, so this never
# oscillates combo<->solo within a boot on its own.
#
# Concurrency: no systemd ordering exists against the sibling
# jasper-fanin-coupling-auto.service or a direct operator CLI run. The reconcile
# CLI serializes every entry verb behind an advisory flock
# (/run/jasper-fanin-coupling.lock) with a bounded wait, and aborts loudly
# (unit lands `failed`, doctor-visible) on contention.
[Unit]
Description=USB combo runtime-fallback health check (disarm to aloop bridge on sustained capture break)
Documentation=file:/opt/jasper/docs/HANDOFF-usbsink.md
# Order after the audio graph + control so a disarm transition (which bounces
# fan-in/usbsink through the restart broker and reloads CamillaDSP) lands on live
# daemons. Wants (not Requires) so a degraded box still ticks and no-ops cleanly.
After=jasper-fanin.service jasper-outputd.service jasper-camilla.service jasper-control.service
Wants=jasper-fanin.service jasper-control.service

[Service]
Type=oneshot
# Bounded: one STATUS socket read plus, only on a sustained-broken transition, the
# ordered combo disarm (fan-in + usbsink restart via the broker + a CamillaDSP
# reconcile). A healthy tick is a single socket read.
TimeoutStartSec=60
# Same env chain the reconciler CLI hydrates internally (load_env_files); listed
# here too so the unit environment carries the combo state for logging context.
EnvironmentFile=-/etc/jasper/jasper.env
EnvironmentFile=-/var/lib/jasper/fanin.env
ExecStart=/opt/jasper/.venv/bin/jasper-fanin-coupling-reconcile --health --reason systemd
StandardOutput=journal
StandardError=journal
