# deploy/systemd/jasper-usbgadget-forensics.path
[Unit]
Description=Start Jasper USB gadget forensics when enabled

[Path]
PathExists=/var/lib/jasper/usb_gadget_forensics.env
Unit=jasper-usbgadget-forensics.service

[Install]
WantedBy=multi-user.target
