# deploy/systemd/jasper-voice.service
[Unit]
Description=Jasper voice daemon (wake word + real-time LLM + tools)
After=jasper-camilla.service network-online.target mpd.service
# Wants= (not Requires=) so a camilla restart blip doesn't cascade-stop
# the voice daemon. The daemon is hardened to keep wake-word + cue
# playback working through a 2 s camilla outage (camilla calls go
# through best_effort=True paths in jasper/camilla.py). After= still
# orders us behind camilla on cold boot.
Wants=network-online.target jasper-camilla.service

[Service]
Type=simple
EnvironmentFile=/etc/jasper/jasper.env
# Spotify CLIENT_ID/SECRET written by the setup wizard at
# https://jts.local/spotify. Optional — leading `-` lets the
# daemon boot fine when this file doesn't exist yet (the
# spotify_play tool then short-circuits to a "visit
# JASPER_SPOTIFY_SETUP_URL" error which Gemini speaks).
EnvironmentFile=-/var/lib/jasper/spotify_credentials.env
# Voice provider keys (GEMINI / OPENAI / XAI) and active-provider
# selection written by the wizard at https://jts.local/voice.
# Sourced AFTER /etc/jasper/jasper.env so wizard-written values
# override operator-managed defaults — same pattern as the
# Spotify creds file above. Optional via the leading `-`.
EnvironmentFile=-/var/lib/jasper/voice_provider.env
ExecStart=/opt/jasper/.venv/bin/jasper-voice
WorkingDirectory=/opt/jasper
Restart=on-failure
RestartSec=2
StateDirectory=jasper
StateDirectoryMode=0750
# RuntimeDirectory creates /run/jasper at service start (auto-cleaned
# at stop). voice_daemon binds /run/jasper/voice.sock there for
# jasper-control to send dial hold-to-talk triggers via.
RuntimeDirectory=jasper
RuntimeDirectoryMode=0750

[Install]
WantedBy=multi-user.target
