# deploy/systemd/jasper-aec-reconcile.service
[Unit]
Description=Reconcile microphone and AEC input-profile state
Documentation=file:/opt/jasper/docs/HANDOFF-resilience.md
After=sound.target jasper-camilla.service
Wants=sound.target jasper-camilla.service

[Service]
Type=oneshot
TimeoutStartSec=60
EnvironmentFile=-/etc/jasper/jasper.env
EnvironmentFile=-/var/lib/jasper/aec_mode.env
ExecStart=/usr/local/sbin/jasper-aec-reconcile --reason systemd
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
