# deploy/systemd/jasper-dac-init.service
[Unit]
Description=Pin Apple USB-C dongle Headphone control to max (fixed analog ceiling)
# The dongle's "Headphone" amixer control is an analog attenuator
# in the DAC's output stage. Leaving it below 100% throws away
# headroom we'd want available to the digital chain (CamillaDSP
# main_volume) — the user-perceived volume should be controlled
# entirely in software. Pin it at 100% on every boot so a stray
# alsactl restore or external setting can't quietly cap us.
#
# Runs once at boot, BEFORE jasper-camilla so the speaker output
# path is at full hardware gain by the time camilla opens it.
# Idempotent and runtime-gated: on non-Apple-DAC systems, or when
# the dongle is temporarily absent, the helper exits cleanly.
After=sound.target
Before=jasper-camilla.service shairport-sync.service librespot.service
Requires=sound.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/jasper-dac-init __APPLE_DONGLE_CARD__ Headphone
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
