# deploy/systemd/jasper-xvf-firmware-update.service
[Unit]
Description=Download and update XVF3800 microphone firmware
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/jasper/jasper.env
ExecStart=/opt/jasper/.venv/bin/jasper-xvf-firmware-update
# The updater refuses to enter DFU after 150s of pre-flash work, then can spend
# 271s on its longest handled path: stop services (30s) + DFU flash (120s) +
# re-enumeration/poll (31s) + reconcile (45s) + recovery reconcile (45s).
# Retain at least 30s of process/scheduling headroom so systemd never kills
# dfu-util during an update that remains within the updater's own deadlines.
TimeoutStartSec=480

# Runs as root because DFU, service stop/start, and reconciler handoff are
# privileged hardware operations. This unit is only started by explicit user
# action through jasper-control's scoped polkit allowlist.
User=root
Group=root
UMask=0022

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/jasper
RestrictNamespaces=true
LockPersonality=true
SystemCallArchitectures=native
