# deploy/systemd/jasper-fanin.service
[Unit]
Description=Jasper renderer fan-in (per-substream → summed music reference)
Documentation=file:///opt/jasper/docs/HANDOFF-fan-in-daemon.md
After=sound.target
Before=jasper-camilla.service shairport-sync.service librespot.service bluealsa-aplay.service jasper-usbsink.service jasper-aec-bridge.service

# StartLimit tuned for the audio path. Tighter burst than the
# default-4/300 sidecar pattern because fan-in is hot-path
# infrastructure (CamillaDSP + AEC bridge depend on its output);
# but generous enough to absorb transient hardware blips
# (USB DAC de-enumeration, USB suspend resume). After 5 failures
# in 5 minutes we escalate to a clean system reboot (Tier 5.1) —
# matches the resilience-ladder pattern for jasper-camilla /
# jasper-aec-bridge / jasper-voice / jasper-control.
StartLimitIntervalSec=300
StartLimitBurst=5
StartLimitAction=reboot

[Service]
# Type=notify + WatchdogSec=30s is the JTS-standard Tier 1+2
# resilience contract. The daemon sends READY=1 after its ALSA PCMs
# are initialized and the STATUS server thread is spawned, then pings
# WATCHDOG=1 every ~10 s, gated on a progress sentinel updated by the
# work loop after every successful frame. A wedged work loop stops
# bumping; the heartbeat thread stops pinging; systemd reaps us within
# WatchdogSec.
Type=notify
WatchdogSec=30s

# 5 s, not the systemd default 90 s. Critical to escalate to
# SIGKILL fast — the 2026-05-11 snd-aloop wedge taught the project
# that a daemon with blocked I/O can sit on SIGTERM for 90 s and
# corrupt kernel-side ALSA state by the time SIGKILL fires. 5 s
# beats that race.
TimeoutStopSec=5s

# Restart=on-failure covers all non-zero-exit paths including
# WatchdogSec timeout. RestartSec=5 backs off enough for a
# transient hardware blip (e.g., USB de-enumeration) to settle
# before retry.
Restart=on-failure
RestartSec=5

# OOM ladder slot. Between Camilla (-900, silence-critical) and
# AEC bridge (-700, capture-critical). The fan-in daemon is the
# upstream source of the music signal both consume; killing it
# preferentially over Camilla makes sense (Camilla can emit silence
# briefly; losing Camilla itself means total silence).
# See docs/HANDOFF-resilience.md "OOM ladder" section.
OOMScoreAdjust=-800

# Stage 2 audio protection: pages never paged to zram. Audio jitter
# from zram decompression latency is the dominant risk on a 1 GB
# Pi 5 under memory pressure; this slice membership shields the
# work loop. See docs/HANDOFF-resilience.md "Stage 2 audio-protection
# subset" section. The slice itself is defined in
# deploy/systemd/jts-audio.slice.
Slice=jts-audio.slice

# Real-time scheduling. SCHED_FIFO priority 30 + mlockall.
# OSPERT 2024 measured Pi 5 stock-kernel worst-case scheduling
# latency at 36.8 ms under heavy stress; SCHED_FIFO + mlockall
# materially reduces this. PREEMPT_RT would reduce further (to
# ~125 µs) and is the documented escalation path if the 72-hour
# soak under realistic stress shows xruns.
LimitMEMLOCK=infinity
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=30

# LimitRTTIME bounds RT-thread runaway (audio-latency foundation G4). A FIFO
# thread that spins without blocking can starve PID 1 on the 1 GB Pi and trip
# the hardware watchdog into a full reboot; SIGXCPU at 200 ms turns that into a
# single crashed daemon that Restart=on-failure revives instead. Mirrors
# PipeWire RTKit's RTTimeUSec. Carried by every FIFO unit (camilla/fanin/outputd).
LimitRTTIME=200000

# Hardening — matches the conventions of jasper-voice,
# jasper-outputd. Defense in depth.
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/var/lib/jasper /run/jasper-fanin

# Runtime dir auto-created on start. systemd cleans it up on stop, so the UDS
# socket never leaks across daemon lifecycle.
# WS1 Phase 3b: fanin stays root (final-output path), but joins the `jasper`
# group + sets UMask=0007 so its tts.sock / control.sock are created
# root:jasper 0770 (a UNIX socket needs WRITE permission to connect(), so the
# old umask-derived 0755 blocked the now-non-root jasper-voice/jasper-mux — only
# root could connect). RuntimeDirectoryMode=0750 (root:jasper) lets the group
# traverse to the sockets.
Group=jasper
UMask=0007
RuntimeDirectory=jasper-fanin
RuntimeDirectoryMode=0750

# Configuration files — chain matches voice / aec / control daemons.
# /etc/jasper/jasper.env is the operator's system-wide config.
# The wizard/reconciler-owned /var/lib/jasper/fanin.env is loaded LAST
# (after the Environment= defaults below) so its overrides win — systemd
# applies a later EnvironmentFile= over an earlier Environment=. This is
# load-bearing for jasper.fanin.buffer_reconcile: it writes the adaptive
# JASPER_FANIN_OUTPUT_BUFFER_FRAMES override into that file, which MUST
# beat the Environment="...=1024" default below. Keep the wizard file's
# EnvironmentFile= line AFTER the Environment= lines (see end of section).
EnvironmentFile=/etc/jasper/jasper.env

# Input ALSA buffer sizing — must absorb the WiFi delivery burst
# pattern. Output buffering is split out below and intentionally kept
# smaller so the burst absorber does not become hidden downstream
# latency between fanin and CamillaDSP/AEC.
#
# 802.11 A-MPDU aggregation batches multiple AirPlay RTP packets into a
# single radio transmission; the on-wire pattern observed via tcpdump
# (2026-05-26) is **bursts of ~4 packets every ~30 ms** instead of the
# nominal 1 packet every 8 ms at 44.1 kHz. shairport faithfully writes
# those bursts into its snd-aloop substream, so fanin's per-input ALSA
# ring sees the same bursty fill pattern.
#
# With the upstream default of buffer_frames=1024 (~21 ms at 48 kHz),
# the ring would overrun before the next burst arrived — measured at
# ~2 xruns/min on real hardware (every xrun = one period of injected
# silence = audible click every 30-60 s). With 4096 (~85 ms), the ring
# absorbs the worst-case ~40 ms burst gap comfortably. **This matches
# the buffer size of the userspace dmix layer (PR #214) that fanin
# replaces** — that dmix was, accidentally, the WiFi-burst absorption
# layer. fanin must provide the same absorption capacity to be
# regression-free against dmix.
#
# Trade-off: up to ~85 ms of input-side queue capacity, same as the
# dmix layer had. This is headroom, not a fixed delay added to every
# frame: actual queued audio depends on the writer/reader fill level,
# and shairport can observe snd-aloop delay on its private lane.
#
# The output side is intentionally smaller than the input burst absorber.
# 1024 frames (~21 ms) is validated on the low-latency Apple-DAC path with
# CamillaDSP chunksize=256/target_level=1536. That older 3072-frame default
# was needed when the downstream Camilla path used a 1024-frame chunk; with
# the DAC floor applied, it is latency without useful headroom.
#
# Read at runtime by jasper_fanin::config::Config::from_env().
# Override via /var/lib/jasper/fanin.env if you need to retune
# (e.g., a future PR with verified-tighter WiFi delivery).
Environment="JASPER_FANIN_INPUT_BUFFER_FRAMES=4096"
Environment="JASPER_FANIN_OUTPUT_BUFFER_FRAMES=1024"

# Wizard/reconciler-owned overrides — loaded LAST so they win over the
# Environment= defaults above (later EnvironmentFile= beats earlier
# Environment=). jasper.fanin.buffer_reconcile is the single writer of the
# adaptive JASPER_FANIN_OUTPUT_BUFFER_FRAMES override here. Optional (leading
# `-`): absent on a fresh box, where the 1024 default above is live.
EnvironmentFile=-/var/lib/jasper/fanin.env

ExecStart=/opt/jasper/bin/jasper-fanin

# Combo-mode host-slaved USB clock (C6) — belt-and-braces for the
# pitch-neutrality safety invariant. In combo mode
# (JASPER_FANIN_USB_DIRECT=enabled + JASPER_FANIN_HOST_CLOCK=enabled) the
# `fanin-host-clock` thread owns the gadget's "Capture Pitch 1000000" ctl and
# resets it to neutral on every clean exit path — but SIGKILL / OOM-kill /
# watchdog abort skip that in-process cleanup, and a host must never stay
# slaved to a command from a daemon that is no longer running.
#
# The owner gate is LOAD-BEARING, not cosmetic: fan-in restarts on every deploy,
# and an UNCONDITIONAL neutralize here would, on a solo-mode box (usbsink owns the
# clock in aloop mode), overwrite an active usbsink L0 command with 1000000 behind
# usbsink's back — desyncing its write-suppression epsilon gate (it believes the
# last written value is its own command and won't rewrite until >10 ppm drift).
# So the belt fires ONLY when THIS daemon is the configured clock owner.
#
# Ownership requires BOTH flags, exactly matching the daemon: fan-in resolves
# `host_clock_enabled_effective = host_clock_enabled && !usb_direct_off` and
# never writes a non-neutral pitch unless BOTH JASPER_FANIN_HOST_CLOCK and
# JASPER_FANIN_USB_DIRECT are set (with only HOST_CLOCK set it logs
# `noop reason=usb_direct_off` and issues ZERO ctl writes — review F2). Gating on
# HOST_CLOCK alone fires the belt on a box that part-rolled-back combo (unset
# USB_DIRECT, left HOST_CLOCK=enabled) while solo usbsink's own DLL is the live
# ctl writer — neutralizing usbsink's command on every fan-in stop, i.e. every
# deploy. Requiring USB_DIRECT too loses zero coverage: fan-in owns the ctl only
# when both flags were set in the same env snapshot ExecStopPost sees.
#
# Both checks are CASE-INSENSITIVE to match the config parser, which arms each
# flag on `eq_ignore_ascii_case("enabled")` (jasper_fanin::config, mirroring
# usbsink's literal idiom), and the `.env.example` prose that advertises
# "case-insensitive" (review F3). An exact-match `= "enabled"` here would leave a
# box armed with e.g. `JASPER_FANIN_HOST_CLOCK=Enabled` running the servo with a
# DEAD belt — the precise non-neutral-after-SIGKILL scenario C6 exists to cover.
# `tr` lowercases each value before the compare.
#
# The card is derived from JASPER_FANIN_USB_DIRECT_DEVICE (default hw:UAC2Gadget)
# so a device override redirects the belt to the right card (review N6 — no
# hardcoded card drift; the in-daemon actuator derives the same card via
# `ctl_card_from_capture`). Mirrors jasper-usbsink.service's owner-gated
# ExecStopPost so both writers target the same element the same way.
#
# The gate + card-derive + neutralize live in the shipped helper
# jasper-fanin-pitch-neutralize (deploy/bin/, installed to /usr/local/sbin by
# install.sh). It was moved OUT of an inline `sh -c` here (defect E, 2026-07-05):
# the inline form derived the card with `card="${card%%,*}"`, and systemd
# consumes `%%` as a specifier escape — logging "Invalid environment variable
# name evaluates to an empty string: card%,*" on every start AND running the
# SHORTEST-match `${card%,*}` instead of the longest-match at runtime (worked by
# luck for the current bare-card device names). No `%` gymnastics in a real
# script. Behavior is identical: fires ONLY when THIS daemon is the configured
# clock owner (BOTH JASPER_FANIN_HOST_CLOCK and JASPER_FANIN_USB_DIRECT set to
# `enabled`, case-insensitive) → amixer reset to neutral 1000000. `-` ignores a
# missing card (combo off, gadget torn down); the helper also exits 0 on the
# not-owner / absent-card paths so the journal stays clean.
ExecStopPost=-/usr/local/sbin/jasper-fanin-pitch-neutralize

[Install]
WantedBy=multi-user.target
