# deploy/systemd/jasper-usbsink-volume.service
[Unit]
Description=Jasper USB sink host-volume observer
# Non-real-time helper for USB Audio Input. jasper-fanin owns the audio data
# plane; this Python process only polls the UAC2 gadget mixer and reports host
# slider changes to jasper-control.
#
# It tracks the process-free readiness marker's lifecycle: the UAC2 mixer exists
# only while the uac2 function is composed. The composite gadget outlives audio
# because it also carries the always-on network.
After=jasper-usbsink.service jasper-control.service
Requires=jasper-usbsink.service
PartOf=jasper-usbsink.service

[Service]
Type=simple
TimeoutStartSec=2s
TimeoutStopSec=5s
EnvironmentFile=-/etc/jasper/jasper.env
UnsetEnvironment=LD_PRELOAD LD_LIBRARY_PATH LD_AUDIT GLIBC_TUNABLES
ExecCondition=+/usr/bin/env -i PATH=/opt/jasper/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin /opt/jasper/.venv/bin/jasper-local-source-allowed --source usbsink
ExecStart=/opt/jasper/.venv/bin/jasper-usbsink-volume
Restart=on-failure
RestartSec=2s
Nice=5
OOMScoreAdjust=100
MemoryHigh=32M
MemoryMax=64M

[Install]
WantedBy=multi-user.target
