# deploy/systemd/jasper-usbsink-init.service
[Unit]
Description=Jasper USB sink — ConfigFS UAC2 gadget setup
# Runs once when jasper-usbsink starts. Loads libcomposite, writes the
# ConfigFS gadget descriptor under /sys/kernel/config/usb_gadget/
# jts-usb-audio, and waits for /proc/asound/UAC2Gadget to appear.
#
# The dtoverlay `dwc2,dr_mode=peripheral` is added to
# /boot/firmware/config.txt by install.sh's set_usb_gadget_mode()
# (requires a reboot to take effect). That puts the BCM2712 OTG
# controller into peripheral mode at boot — but we don't load
# libcomposite or create the gadget descriptor at boot. Both happen
# only when this unit starts, so RAM stays at baseline (~50 KB
# dwc2 module) when the USB sink feature is disabled.
#
# Stop semantics: ExecStop tears down the ConfigFS tree, ExecStopPost
# rmmods libcomposite + u_audio. After a clean stop the only resident
# trace is dwc2 itself.
After=systemd-modules-load.service sys-kernel-config.mount
Before=jasper-usbsink.service
ConditionPathExists=/sys/kernel/config

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/var/lib/jasper/speaker_name.env
# modprobe libcomposite drags in u_audio + configfs deps. Failing
# modprobe means the kernel is missing CONFIG_USB_GADGET / CONFIG_USB_
# CONFIGFS — surface that as a unit failure rather than letting the
# next step fail mysteriously.
ExecStartPre=/sbin/modprobe libcomposite
# Make the host-visible device name track the Speaker Name by patching
# the kernel's hardcoded UAC2 AudioStreaming string (configfs can't set
# it on 6.12). Best-effort (leading `-`): a failure here only means the
# host shows the default "Playback Inactive" label — it must never block
# USB audio. Runs before gadget-up so the patched module autoloads when
# the function is created. See docs/HANDOFF-usbsink.md "Device name".
ExecStartPre=-/usr/local/sbin/jasper-usbsink-name-patch
ExecStart=/usr/local/sbin/jasper-usbsink-gadget-up
ExecStartPost=/usr/local/sbin/jasper-usbsink-wait-card 30
ExecStop=/usr/local/sbin/jasper-usbsink-gadget-down
# Best-effort rmmod — if u_audio or libcomposite is still in use
# (e.g. a stale gadget descriptor or a stuck capture), the rmmod
# fails harmlessly. The leading `-` tells systemd to ignore non-zero
# exit. Either way the gadget descriptor is gone, so the host stops
# seeing JTS as an audio device.
ExecStopPost=-/sbin/rmmod u_audio
ExecStopPost=-/sbin/rmmod libcomposite
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
