# deploy/systemd/jasper-aec-bridge.service
[Unit]
Description=Software AEC bridge — WebRTC AEC3 cancels speaker echo from XVF raw mic
# Captures from pcm.jasper_capture (host-clocked dsnoop on the
# renderer→camilla loopback, the AEC reference) AND from
# hw:CARD=Array,DEV=0 channel 2 (raw mic 0, exposed by 6-ch
# firmware variant 2.0.8). Outputs the AEC'd mono mic over UDP to
# 127.0.0.1:JASPER_AEC_UDP_PORT (default 9876) for jasper-voice to
# receive via UdpMicCapture (JASPER_MIC_DEVICE=udp:9876).
#
# Historical note: this used to write to a second snd-aloop card
# (LoopbackAEC) that jasper-voice read via PortAudio. Replaced
# with UDP in PR 2 of the resilience-ladder series — snd-aloop's
# kernel-side loopback_cable wedged on consumer SIGKILL and
# required a reboot to clear. UDP has no such state.
#
# Engine: WebRTC AEC3 via the jasper_aec3 pybind11 binding around
# Trixie's libwebrtc-audio-processing-1 v1.3-3. Tuning knobs
# (JASPER_AEC_REF_GAIN_DB, JASPER_AEC_MIC_GAIN_DB, JASPER_AEC_AGC2)
# come from /etc/jasper/jasper.env — see docs/HANDOFF-aec.md
# "Tuning findings" for the chosen production values.
#
# REPLACES the previous CamillaDSP-based bridge that was trying to
# feed the XVF chip's USB-IN as far-end reference. The chip's AEC
# turned out to be architecturally incompatible with our external-
# DAC topology — see jasper/cli/aec_bridge.py header for the
# detailed analysis.
After=jasper-camilla.service jasper-aec-init.service sound.target
Requires=jasper-camilla.service
# `Wants=` so init runs whenever bridge starts (init is a one-shot
# that applies XVF tuning; bridge depends on chip being initialised).
# Soft dep — bridge still tries to start even if init failed; AEC
# quality degrades but the daemon is operational.
Wants=jasper-aec-init.service
# `PartOf=` propagates camilla's stop / restart to this unit. Without
# it, a camilla restart (e.g. during deploy) deactivates the bridge
# via `Requires=` cascade — but `Restart=on-failure` doesn't trigger
# on a clean exit (code 0 from cascade-stop), so the bridge stays
# dead silently. Mic ends up listening to a (now-unfed) loopback,
# wake events stop firing, speaker appears alive but is deaf. Hit
# in production on 2026-05-10.
PartOf=jasper-camilla.service

[Service]
# Tier 1+2 of the resilience ladder. `Type=notify` makes the daemon
# emit `READY=1` (via sdnotify, see jasper/watchdog.py) before
# systemd considers it active. `WatchdogSec=30s` requires the daemon
# to pat the watchdog every < 30 s; the daemon's heartbeat thread
# only pats if its work loop has made progress recently (sentinel
# pattern). If the loop wedges (PortAudio blocked in a syscall
# holding the GIL — the 2026-05-11 incident), `WATCHDOG=1` stops
# firing, this timer expires, and `Restart=on-watchdog` revives us
# with a fresh process. Critically, this happens BEFORE the SIGTERM
# timeout / SIGKILL that would otherwise corrupt snd-aloop's
# kernel-side loopback_cable state and require a reboot.
#
# `TimeoutStopSec=5s` is the matching shorter SIGTERM grace period
# (down from systemd's default 90 s). If the daemon can't shut down
# in 5 s, it's wedged and we want to get to SIGKILL+restart fast
# rather than accumulating wedge time.
Type=notify
WatchdogSec=30s
TimeoutStopSec=5s
# Pull JASPER_AEC_REF_GAIN_DB, JASPER_AEC_MIC_GAIN_DB, JASPER_AEC_AGC2
# from the same env file jasper-voice uses. Leading `-` makes systemd
# skip the file silently if it's absent (e.g. fresh install before
# install.sh has run); the bridge defaults are conservative
# (ref_gain=0 dB, mic_gain=0 dB, agc2=off) so a missing env file is
# safe — but the production tuning values won't be applied.
EnvironmentFile=-/etc/jasper/jasper.env
ExecStart=/opt/jasper/.venv/bin/jasper-aec-bridge
# `Restart=on-failure` covers non-zero exit + signals + watchdog
# timeouts (per the table in systemd.service(5)). The previous
# `on-watchdog` setting only restarted on watchdog timeout, so the
# explicit `BridgeStalled` exit (code 1, raised by _aec_loop when
# the mic queue stays empty for JASPER_AEC_STALL_RESTART_SEC) was
# NOT being restarted — the service was just parked as failed, and
# the wake-word path was offline until manual intervention. Hit in
# production 2026-05-11.
Restart=on-failure
RestartSec=2
Nice=-10
IOSchedulingClass=realtime
IOSchedulingPriority=4
LimitRTPRIO=99
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
