# deploy/systemd/jasper-usbsink.service
[Unit]
Description=Jasper USB sink — audio bridge from UAC2 gadget into the music chain
# Bridges the host's USB audio output (captured via the ConfigFS UAC2
# gadget set up by jasper-usbsink-init) into the USB-in private
# snd-aloop lane. jasper-fanin sums it with AirPlay/Spotify/BT before
# CamillaDSP, so it is just another renderer from the music chain's
# perspective.
#
# Ships DISABLED. The /sources/ wizard toggles this intent unit. Starting
# it pulls jasper-usbsink-init in via Requires=; stopping this unit alone does
# not tear down the host-visible gadget, so off/parking paths explicitly stop
# jasper-usbsink-init too. When both are off, libcomposite isn't loaded and
# the daemon doesn't run — total resident cost is the ~50 KB dwc2 kernel
# module from the dtoverlay, nothing in user space.
#
# See docs/HANDOFF-usbsink.md for the full design.
After=jasper-usbsink-init.service jasper-fanin.service jasper-camilla.service sound.target
Requires=jasper-usbsink-init.service
Wants=jasper-fanin.service
# PartOf= propagates a stop/restart of the init service to the daemon.
# If init is stopped (toggle off in /sources/), this stops too. If init
# crashes and systemd restarts it, the daemon restarts in sync.
PartOf=jasper-usbsink-init.service

[Service]
# Tier 1 + 2 of the resilience ladder (see docs/HANDOFF-resilience.md).
# Type=notify with WatchdogSec=15s: the daemon's diagnostic loop pats
# the watchdog only when the playback/output callback has shown
# forward progress recently (sentinel pattern via
# jasper.watchdog.Heartbeat). Capture may be idle for hours when a host
# is paused, unplugged, or using another source; output silence is still
# useful work. If the output callback wedges, pats stop and systemd
# kills + restarts within WatchdogSec.
#
# TimeoutStopSec=5s matches the AEC bridge's shorter SIGTERM grace
# period. We'd rather hit SIGKILL+restart fast than accumulate wedge
# time.
Type=notify
WatchdogSec=15s
TimeoutStopSec=5s
# Pull operator-supplied tunables (capture/playback device overrides,
# log level) from the same env file the rest of the stack reads.
# Leading `-` makes systemd skip silently if the file isn't there.
EnvironmentFile=-/etc/jasper/jasper.env
EnvironmentFile=-/var/lib/jasper/usbsink.env
Environment="JASPER_USBSINK_PLAYBACK_DEVICE=usbsink_substream"
ExecCondition=/opt/jasper/.venv/bin/jasper-local-source-allowed
ExecStart=/opt/jasper/.venv/bin/jasper-usbsink
Restart=on-failure
RestartSec=2s
Nice=-10
# Join the zram-shielded audio tier. usbsink is a real-time music source whose
# capture/playback run on fixed-deadline (~10 ms) PortAudio callbacks; off this
# slice its pages can swap to zram, and zram-decompression jitter makes a callback
# miss its deadline in bursts -> snd-aloop xruns -> the bounded queue overflows
# (dropped_full). jts-audio.slice carries MemorySwapMax=0 (the swap shield) — a
# MEMORY policy, not a CPU cap, so it respects the no-CPU-caps rule. usbsink was the
# only music-path daemon left OUTSIDE this slice (diagnosed 2026-06-28 as the dominant
# cause of the bursty USB-drop tail). OOMScoreAdjust sits below the output chain
# (outputd -950 / camilla -900 / fanin -800): a killed source restarts and only USB
# stops, whereas a killed output stops all audio.
# Do NOT add CPUSchedulingPolicy=fifo to this unit — that unit-level RT change is what
# SIGKILL-crash-looped jasper-aec-bridge on 2026-06-27; RT priority, if ever needed,
# must be elected in-thread AFTER start() completes.
Slice=jts-audio.slice
OOMScoreAdjust=-600
# Memory caps. ~22 MB Pss is the expected resident size when running;
# 64 MiB hard cap catches any future memory leak before OOM-killer at
# the cgroup boundary. MemoryHigh is the throttle threshold.
MemoryHigh=48M
MemoryMax=64M
# Runtime dir for the state file (Phase 3+ will write
# /run/jasper-usbsink/state.json here). Owned by the service.
RuntimeDirectory=jasper-usbsink
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
