# deploy/systemd/jasper-voice.service
[Unit]
Description=Jasper voice daemon (wake word + real-time LLM + tools)
After=jasper-camilla.service network-online.target
# Wants= (not Requires=) so a camilla restart blip doesn't cascade-stop
# the voice daemon. The daemon is hardened to keep wake-word + cue
# playback working through a 2 s camilla outage (camilla calls go
# through best_effort=True paths in jasper/camilla.py). After= still
# orders us behind camilla on cold boot.
Wants=network-online.target jasper-camilla.service

# Cap the crash-restart spiral. On 2026-05-10 we found NRestarts=9032
# accumulated silently because the daemon kept hitting transient
# audio-device unavailability (Apple dongle de-enumerated mid-run),
# `Restart=on-failure` + `RestartSec=2` re-tried every 2 s forever,
# and nothing told the operator. With these limits, 20 failures in
# 5 min puts the unit in `failed` state — visible to `systemctl
# status` and `jasper-doctor` — instead of an invisible loop.
StartLimitBurst=20
StartLimitIntervalSec=300

[Service]
# Tier 1+2 of the resilience ladder. `Type=notify` makes the daemon
# emit `READY=1` (via sdnotify, see jasper/watchdog.py) before
# systemd considers it active. `WatchdogSec=30s` requires the
# daemon to pat the watchdog every < 30 s; the heartbeat thread
# only pats if the WakeLoop's mic-frame iterator has advanced
# recently (sentinel pattern). If the async loop wedges or mic
# capture dies (PortAudio InputStream death after USB hiccup,
# blocked I/O, deadlock), `WATCHDOG=1` stops firing and systemd
# revives us with a fresh process. `TimeoutStopSec=5s` cuts the
# default 90 s SIGTERM grace — a wedged daemon should get to
# restart fast, not accumulate wedge time.
Type=notify
WatchdogSec=30s
TimeoutStopSec=5s
EnvironmentFile=/etc/jasper/jasper.env
# Spotify CLIENT_ID/SECRET written by the setup wizard at
# https://jts.local/spotify. Optional — leading `-` lets the
# daemon boot fine when this file doesn't exist yet (the
# spotify_play tool then short-circuits to a "visit
# JASPER_SPOTIFY_SETUP_URL" error which Gemini speaks).
EnvironmentFile=-/var/lib/jasper/spotify_credentials.env
# Voice provider keys (GEMINI / OPENAI / XAI) and active-provider
# selection written by the wizard at https://jts.local/voice.
# Sourced AFTER /etc/jasper/jasper.env so wizard-written values
# override operator-managed defaults — same pattern as the
# Spotify creds file above. Optional via the leading `-`.
EnvironmentFile=-/var/lib/jasper/voice_provider.env
# Google OAuth client ID/SECRET written by the wizard at
# http://jts.local/google. Calendar + Gmail tools register only
# when this file's CLIENT_ID/SECRET are set AND at least one
# household member has linked an account. Optional via leading `-`.
EnvironmentFile=-/var/lib/jasper/google_credentials.env
# Wake-word selection (JASPER_WAKE_MODEL) written by the picker at
# http://jts.local/wake/, or seeded by install.sh on a fresh install
# to point at /var/lib/jasper/wake/jarvis_v2.onnx. Sourced AFTER
# /etc/jasper/jasper.env so wizard-written values win. Optional via
# the leading `-` — daemon falls back to bundled "hey_jarvis" when
# the file is missing. See jasper/wake_models.py for the registry.
EnvironmentFile=-/var/lib/jasper/wake_model.env
# Transit settings (subway station, bus stop, MTA BusTime key, etc.)
# written by the wizard at http://jts.local/transit/. Subway and bus
# tools register only when the relevant env vars are populated, so a
# fresh install with no transit configured leaves the model with no
# transit tool surface — by design. Same override pattern as the
# files above. Optional via the leading `-`.
EnvironmentFile=-/var/lib/jasper/transit.env
ExecStart=/opt/jasper/.venv/bin/jasper-voice
WorkingDirectory=/opt/jasper
# `on-watchdog` covers BOTH clean failure (non-zero exit, segfault)
# AND watchdog timeouts. Strictly broader than `on-failure`.
Restart=on-watchdog
# 5 s backoff (was 2 s). Audio-device transients (dongle enumerate
# blip, ALSA card reorder) typically settle within a few seconds; a
# longer backoff means fewer wasted restart cycles before the device
# stabilises, and pairs with StartLimitBurst above to cap the spiral.
RestartSec=5
StateDirectory=jasper
StateDirectoryMode=0750
# RuntimeDirectory creates /run/jasper at service start (auto-cleaned
# at stop). voice_daemon binds /run/jasper/voice.sock there for
# jasper-control to send dial hold-to-talk triggers via.
RuntimeDirectory=jasper
RuntimeDirectoryMode=0750

# Hardening — defense in depth around the highest-attack-surface
# process on the speaker (parses microphone audio, third-party LLM
# responses, Spotify/Google API responses). Runs as root because it
# needs realtime audio + multiple /dev/snd/* + ALSA config in /root,
# but these flags limit damage if anything ever triggers an
# exploitable bug. Mirrors jasper-web.service. ProtectHome=read-only
# (not =true) so ALSA can still read /root/.asoundrc.
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/var/lib/jasper

[Install]
WantedBy=multi-user.target
