# deploy/systemd/jasper-outputd.service
[Unit]
Description=Jasper final output owner
Documentation=file:///opt/jasper/docs/HANDOFF-speaker-output-reference.md
After=sound.target jasper-fanin.service jasper-camilla.service

StartLimitIntervalSec=300
StartLimitBurst=5
# outputd owns the final DAC boundary. If it cannot stay up after the
# normal restart burst, reboot cleanly rather than leaving the speaker
# silently pinned without a final-output owner.
StartLimitAction=reboot

[Service]
Type=notify
WatchdogSec=30s
TimeoutStopSec=5s
Restart=on-failure
RestartSec=5
# CONFIG-validation failures exit 78 (EX_CONFIG) and must PARK the unit
# failed (visible on /state + doctor) instead of restart-looping:
# restarting cannot fix bad config, and on this unit the loop escalates
# to StartLimitAction=reboot. Measured incident (jts3, 2026-06-11): a
# grouping lane env + a lab content-bridge retune layered into a
# guard-rejected combination; outputd crash-looped into THREE Pi reboots
# before the T5.1 boot-loop guard contained it. Recovery after a park:
# fix the env, then `systemctl restart jasper-outputd`.
RestartPreventExitStatus=78

# Final-output owner. This is more important than Camilla for speaker
# silence because it owns the DAC write loop.
OOMScoreAdjust=-950
Slice=jts-audio.slice
LimitMEMLOCK=infinity
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=35

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/var/lib/jasper /run/jasper-outputd
# WS1 Phase 3b: outputd stays root (it owns the DAC write loop), but joins the
# `jasper` group + sets UMask=0007 so its tts.sock (which the now-non-root
# jasper-voice connects to for member-local TTS while bonded) is created
# root:jasper 0770 — connect() needs WRITE permission, so the old umask-derived
# 0755 only let root connect. RuntimeDirectoryMode=0750 (root:jasper) lets the
# group traverse.
Group=jasper
UMask=0007
RuntimeDirectory=jasper-outputd
RuntimeDirectoryMode=0750

EnvironmentFile=/etc/jasper/jasper.env
Environment="JASPER_OUTPUTD_BACKEND=alsa"
Environment="JASPER_OUTPUTD_CONTENT_PCM=outputd_content_capture"
Environment="JASPER_OUTPUTD_DAC_PCM=outputd_dac"
Environment="JASPER_OUTPUTD_PERIOD_FRAMES=1024"
Environment="JASPER_OUTPUTD_CONTENT_BUFFER_FRAMES=4096"
Environment="JASPER_OUTPUTD_DAC_BUFFER_FRAMES=3072"
Environment="JASPER_OUTPUTD_CONTENT_BRIDGE=direct"
Environment="JASPER_OUTPUTD_CONTROL_SOCKET=/run/jasper-outputd/control.sock"
# Optional operator retune file. Keep this after the packaged defaults so
# outputd and jasper-apply-airplay-mode agree when DAC buffer sizing is
# intentionally changed for lab validation.
EnvironmentFile=-/var/lib/jasper/outputd.env
# Wake-corpus test mode can temporarily enable the XVF USB-IN
# reference fanout and outputd's UDP reference tap. Keep this after
# operator retunes so entering/leaving corpus mode is reversible.
EnvironmentFile=-/var/lib/jasper/wake_corpus_bridge.env
# Multi-room round-trip lane (HANDOFF-multiroom.md §2, Increment 5).
# Written ONLY by jasper-grouping-reconcile (the single writer):
# JASPER_OUTPUTD_DAC_CONTENT_FIFO + _CHANNEL wire the Increment 3
# dac_content reader to the member round-trip FIFO + this speaker's
# channel pick. PERSISTENT (not /run) so a bonded speaker boots with
# the lane already configured — no extra restart at boot. Empty values
# (written on disband — disable-clears-stale) read as unset: the
# byte-identical solo DAC loop. Last env layer on purpose: the
# reconciler's derived state must win over any hand-set value.
EnvironmentFile=-/var/lib/jasper/grouping-outputd.env
# If the resolved final-output card disappeared, skip the start before
# the main process opens ALSA. The fake backend deliberately opens no
# card, keeping outputd's sockets/state alive in no-DAC states, so it
# must pass this gate. ExecCondition exit 1 parks the unit inactive
# without consuming Restart=on-failure/StartLimitAction; the
# audio-hardware udev reconciler reset-failed+starts outputd when a
# recognized DAC role returns. Camilla is intentionally not gated here:
# in the outputd cutover topology it writes only to snd-aloop.
ExecCondition=/bin/sh -c 'backend="$${JASPER_OUTPUTD_BACKEND:-alsa}"; card="$${JASPER_AUDIO_DAC_CARD:-}"; id="$${JASPER_AUDIO_DAC_ID:-unknown}"; if [ "$$backend" = "fake" ] || [ -z "$$card" ] || [ -e "/proc/asound/$$card" ]; then exit 0; fi; echo "event=outputd.output_device_gate.park reason=missing_dac output_dac_id=$$id output_dac_card=$$card outputd_backend=$$backend" >&2; exit 1'
# On an outputd crash, refresh the output-hardware env before systemd's
# built-in Restart=on-failure retry. The helper skips normal stops,
# ExecCondition parks, and EX_CONFIG=78; it uses --no-restart so it never
# controls jasper-outputd from inside this unit's stop transaction.
ExecStopPost=-/usr/local/sbin/jasper-outputd-failure-reconcile
ExecStart=/opt/jasper/bin/jasper-outputd

[Install]
WantedBy=multi-user.target
