# deploy/systemd/jasper-voice.service
[Unit]
Description=Jasper voice daemon (wake word + real-time LLM + tools)
After=jasper-fanin.service jasper-camilla.service jasper-outputd.service jasper-aec-bridge.service network-online.target jasper-accessory-reconcile.service
# Wants= (not Requires=) so a camilla restart blip doesn't cascade-stop
# the voice daemon. The daemon is hardened to keep wake-word + cue
# playback working through a 2 s camilla outage (camilla calls go
# through best_effort=True paths in jasper/camilla.py). After= still
# orders us behind camilla on cold boot.
# jasper-aec-bridge is also a soft dependency: on the reconciler-managed
# `udp:9876` mic path it is the UDP producer. If it was clean-stopped by a
# Camilla/outputd restart, starting voice without also starting the bridge makes
# UdpMicCapture bind an unfed socket; the wake loop then stops advancing and
# systemd watchdog-kills voice in a visible restart loop. Keep this as `Wants=`
# (not `Requires=`) so a bridge restart blip does not cascade-stop voice.
Wants=network-online.target jasper-fanin.service jasper-camilla.service jasper-outputd.service jasper-aec-bridge.service jasper-accessory-reconcile.service

# Cap the crash-restart spiral. On 2026-05-10 we found NRestarts=9032
# accumulated silently because the daemon kept hitting transient
# audio-device unavailability (Apple dongle de-enumerated mid-run),
# `Restart=on-failure` + `RestartSec=2` re-tried every 2 s forever,
# and nothing told the operator. With these limits, 20 failures in
# 5 min puts the unit in `failed` state — visible to `systemctl
# status` and `jasper-doctor` — instead of an invisible loop.
StartLimitBurst=20
StartLimitIntervalSec=300

# Tier 4.5 — escalate to system reboot if the per-service restart
# spiral exceeds the existing burst limit (20 in 5min). The 20/300
# values are intentionally tolerant of transient Apple-dongle
# de-enumeration (per the comment block above re: 2026-05-10).
# But if the spiral DOES exceed 20 — which means sustained sickness
# beyond what transient hardware blips can produce — the right
# answer is a clean system reboot, not the kernel hardware
# watchdog at Tier 5. See docs/HANDOFF-tier5-watchdog-liveness.md
# "Option B (T5.1)". `reboot` (not `reboot-force`): we want clean
# shutdown so zram dirty pages sync; on a 1 GB Pi `reboot-force`
# can leave the next boot in EXT4 "orphan cleanup" state.
StartLimitAction=reboot

# Microphone-presence gate (docs/HANDOFF-hotplug-resilience.md "Layer 1").
# jasper-aec-reconcile is the single writer of this PERSISTENT NEGATIVE
# marker: it exists only when the reconciler positively determined there
# is no usable mic. The `!` makes this "run UNLESS the marker exists", so
# it FAILS OPEN — a fresh box, or one where the reconciler never ran,
# starts voice exactly as before. When the marker IS present, PID 1 skips
# the start cleanly (ActiveState=inactive, NOT counted as a start, NOT
# subject to Restart=, NEVER escalates to StartLimitAction=reboot above) —
# the same clean-skip the output owner relies on (jasper-outputd.service
# ExecCondition). This is what stops a no-mic box from crash-looping into
# a reboot; on plug-in, udev → jasper-aec-reconcile removes the marker and
# restart_voice brings us up. Marker path is duplicated as a literal in
# deploy/bin/jasper-aec-reconcile and jasper/voice/input_presence.py; the
# agreement is pinned by tests/test_voice_input_gate.py.
ConditionPathExists=!/var/lib/jasper/voice-input-absent

[Service]
# Tier 1+2 of the resilience ladder. `Type=notify` makes the daemon
# emit `READY=1` (via sdnotify, see jasper/watchdog.py) before
# systemd considers it active. `WatchdogSec=30s` requires the
# daemon to pat the watchdog every < 30 s; the heartbeat thread
# only pats if the WakeLoop's mic-frame iterator has advanced
# recently (sentinel pattern). If the async loop wedges or mic
# capture dies (PortAudio InputStream death after USB hiccup,
# blocked I/O, deadlock), `WATCHDOG=1` stops firing and systemd
# revives us with a fresh process. `TimeoutStopSec=5s` cuts the
# default 90 s SIGTERM grace — a wedged daemon should get to
# restart fast, not accumulate wedge time.
Type=notify
WatchdogSec=30s
TimeoutStopSec=5s
EnvironmentFile=/etc/jasper/jasper.env
EnvironmentFile=-/var/lib/jasper/speaker_name.env
# Spotify CLIENT_ID/OAUTH_MODE written by the setup wizard at
# https://jts.local/spotify. Optional — leading `-` lets the
# daemon boot fine when this file doesn't exist yet (the
# spotify_play tool then short-circuits to a "visit
# JASPER_SPOTIFY_SETUP_URL" error which Gemini speaks).
# WS1 Phase 4b — moved into the group-`jasper-intsecrets` dir alongside
# the Spotify OAuth token cache.
EnvironmentFile=-/var/lib/jasper-intsecrets/spotify_credentials.env
# Voice provider SELECTION (JASPER_VOICE_PROVIDER + per-provider model /
# voice) written by the wizard at https://jts.local/voice. Sourced AFTER
# /etc/jasper/jasper.env so wizard-written values override operator
# defaults. Optional via the leading `-`. The API keys themselves are
# NOT here anymore — see voice_keys.env below (WS1 Phase 4a split).
EnvironmentFile=-/var/lib/jasper/voice_provider.env
# WS1 Phase 4a — the GEMINI / OPENAI / XAI API keys, split out of
# voice_provider.env into the group-`jasper-secrets` dir so only
# jasper-voice + jasper-web can read them (mux/control/input cannot).
# Readable because this unit joins the jasper-secrets supplementary
# group (below). Optional via the leading `-`.
EnvironmentFile=-/var/lib/jasper-secrets/voice_keys.env
# Google OAuth client ID/SECRET written by the wizard at
# http://jts.local/google. Calendar + Gmail tools register only
# when this file's CLIENT_ID/SECRET are set AND at least one
# household member has linked an account. Optional via leading `-`.
# WS1 Phase 4a — moved into the group-`jasper-secrets` dir alongside
# the OAuth token tree (/var/lib/jasper-secrets/google/).
EnvironmentFile=-/var/lib/jasper-secrets/google_credentials.env
# Google Routes API key written by /transit/ for destination ETA /
# directions. Kept separate from OAuth credentials: this is a billable
# Maps API key, not a Google account token.
EnvironmentFile=-/var/lib/jasper-secrets/google_routes.env
# Wake-word selection (JASPER_WAKE_MODEL) written by the picker at
# http://jts.local/wake/, or seeded by install.sh on a fresh install
# to point at /var/lib/jasper/wake/jarvis_v2.onnx. Sourced AFTER
# /etc/jasper/jasper.env so wizard-written values win. Optional via
# the leading `-` — daemon falls back to bundled "hey_jarvis" when
# the file is missing. See jasper/wake_models.py for the registry.
EnvironmentFile=-/var/lib/jasper/wake_model.env
# Weather default location + units written by the wizard at
# http://jts.local/weather/. Sourced AFTER /etc/jasper/jasper.env so
# wizard-written values win. Optional via leading `-`; with no default,
# explicit-location weather questions still work.
EnvironmentFile=-/var/lib/jasper/weather.env
# Transit settings (subway station, bus stop, MTA BusTime key, etc.)
# written by the wizard at http://jts.local/transit/. Subway and bus
# tools register only when the relevant env vars are populated, so a
# fresh install with no transit configured leaves the model with no
# transit tool surface — by design. Same override pattern as the
# files above. Optional via the leading `-`.
EnvironmentFile=-/var/lib/jasper/transit.env
# Home Assistant URL + token (JASPER_HA_URL, JASPER_HA_TOKEN,
# JASPER_HA_AGENT_ID) written by the wizard at
# http://jts.local/ha (PR 2). Optional via leading `-`;
# the home_assistant tool is gated on both URL and token being set,
# so a missing file just leaves smart-home control disabled.
# See docs/HANDOFF-homeassistant.md.
EnvironmentFile=-/var/lib/jasper-intsecrets/home_assistant.env
# Tools the household DISABLED via the /tools/ wizard
# (JASPER_DISABLED_TOOLS). Wizard-owned SSOT at
# /var/lib/jasper/tool_state.env. The daemon reads this set fail-safe
# (jasper.tool_state) to filter the registry — a missing file means
# nothing disabled (every tool on). Sourced here so jasper-doctor and
# env_load see it; optional via the leading `-`.
EnvironmentFile=-/var/lib/jasper/tool_state.env
# Conversation-history capture toggle and retention written by /chat/.
# The daemon also reads this file fresh before each capture write so a
# household toggle takes effect without a restart; sourcing it here keeps
# systemd/env-load parity and exposes non-default DB paths at startup.
EnvironmentFile=-/var/lib/jasper/conversation_history.env
# Multi-device wake arbitration (JASPER_PEERING) written by the Speakers page at
# http://jts.local/rooms/. jasper-voice is a CLIENT of jasper-control's
# peering UDS: Config.peering_enabled reads JASPER_PEERING here to decide
# whether WakeLoop arbitrates each wake with peers or behaves solo. Without
# this line the daemon never saw JASPER_PEERING, so peering stayed off even
# after the wizard enabled it. Sourced AFTER /etc/jasper/jasper.env so
# wizard values win; optional via leading `-` (absent => solo). MUST match
# jasper-control.service, which also sources it — both daemons read it.
EnvironmentFile=-/var/lib/jasper/peering.env
# Optional accessory-backed push-to-talk mic sources. Single writer:
# jasper-accessory-reconcile, which activates entries only when the matching
# hardware profile is paired (for example WiiM Remote 2 → wiim_remote_2=udp:9892).
EnvironmentFile=-/var/lib/jasper/accessory-mics.env
# Assistant audio routes to fan-in's outputd-compatible local socket so
# TTS/cues enter before CamillaDSP, then get the same crossover,
# correction, and protection path as renderer audio. The transport name
# remains `outputd` for Python API compatibility; the socket endpoint is
# the runtime owner.
Environment="JASPER_TTS_TRANSPORT=outputd"
Environment="JASPER_TTS_OUTPUTD_SOCKET=/run/jasper-fanin/tts.sock"
Environment="JASPER_DUCK_TRANSPORT=fanin"
# Grouping-derived TTS playout route (Increment 5 PR-2): the reconciler's
# route matrix either points JASPER_TTS_OUTPUTD_SOCKET at outputd's socket
# for passive non-sub members, omits it so the inline fanin default above
# applies, and/or parks voice with JASPER_GROUPING_VOICE_PARK. Reconciler-owned,
# layered LAST so the derived state wins. See docs/HANDOFF-multiroom.md §2.
EnvironmentFile=-/var/lib/jasper/grouping-voice.env
ExecStart=/opt/jasper/.venv/bin/jasper-voice
WorkingDirectory=/opt/jasper
# `on-failure` restarts on non-zero exit, signal-termination, an
# operation timeout, OR a watchdog timeout — strictly broader than
# `on-watchdog`, which covers only the watchdog. The prior comment
# in this file got the semantics inverted (claimed `on-watchdog`
# covered non-zero exits); under `on-watchdog`, a daemon that exits
# non-zero from its initial-connect path was NOT restarted. That
# directly caused the 2026-05-23 silent-dead outage where the
# voice loop's OpenAI initial-connect retry budget exhausted
# during a WiFi-recovery window. The initial-connect path now
# uses a 10-minute time budget (see jasper/voice/openai_session.py:
# ``DEFAULT_INITIAL_CONNECT_BUDGET_SEC``) and exits non-zero on
# exhaustion so systemd's outer loop (this `Restart=on-failure` +
# the StartLimitBurst above) grants another full budget.
Restart=on-failure
# Two exit codes mean "intentionally idle, not a crash" — park the unit
# cleanly (treated as success, no auto-restart) so neither can consume the
# StartLimitAction=reboot budget above:
#   78 (EX_CONFIG)  — JASPER_VOICE_PROVIDER unset/invalid (first-time setup
#                     state, before the /voice wizard is configured).
#   66 (EX_NOINPUT) — primary microphone could not be opened at startup
#                     (VOICE_MIC_UNAVAILABLE_EXIT). The ConditionPathExists
#                     marker above pre-empts the common no-mic case; this is
#                     the backstop for custom/present-but-unopenable/first-boot
#                     mics. jasper-aec-reconcile (udev) restarts us on plug-in.
#                     See docs/HANDOFF-hotplug-resilience.md "Layer 2".
SuccessExitStatus=66 78
RestartPreventExitStatus=66 78
# 5 s backoff (was 2 s). Audio-device transients (dongle enumerate
# blip, ALSA card reorder) typically settle within a few seconds; a
# longer backoff means fewer wasted restart cycles before the device
# stabilises, and pairs with StartLimitBurst above to cap the spiral.
RestartSec=5
# WS1 Phase 3b — user drop (3b-1). jasper-voice runs as a dedicated non-root
# service user in the shared `jasper` group so it can reach the cross-daemon
# /run sockets and group-shared /var/lib/jasper files; `audio` (supplementary)
# covers ALSA — the XVF mic is the "Array" ALSA card, NOT a raw USB endpoint
# (jasper-aec-init, which stays root, is the only raw-USB user), so no udev
# rule is needed. jasper-voice makes ZERO realtime-scheduling / mlock syscalls
# (it's a Nice daemon; jasper-aec-bridge is the RT process and stays root), so
# it needs no capability — CapabilityBoundingSet= drops them all and
# SystemCallFilter=@system-service blocks the dangerous syscall classes
# (NoNewPrivileges + SystemCallArchitectures=native below make the filter
# enforce). docs/HANDOFF-privilege-separation.md.
User=jasper-voice
Group=jasper
# WS1 Phase 4a — `jasper-secrets` (supplementary) grants read of the
# split-out LLM API keys + Google client secret/token tree in
# /var/lib/jasper-secrets (which mux/control/input are NOT in). WS1 Phase
# 4b — `jasper-intsecrets` grants read/write of HA + Spotify integration
# secrets/caches (input is NOT in it). `audio` covers ALSA.
# docs/HANDOFF-privilege-separation.md "Phase 4".
SupplementaryGroups=audio jasper-secrets jasper-intsecrets
CapabilityBoundingSet=
SystemCallFilter=@system-service
StateDirectory=jasper
# 0770 (was 0750): the shared `jasper` group must write
# /var/lib/jasper/speaker_volume.json, which jasper-mux and jasper-control
# also read/write. atomic_write_text creates a tempfile in this dir + renames,
# so group write on the dir is required; per-file modes still gate the rest.
StateDirectoryMode=0770
# UMask=0007 makes files this daemon CREATES here group-writable (0660) — the
# "per-file modes" the dir-mode note above defers. usage.db, wake-events.sqlite3,
# timers.db and speaker_volume.json are written by jasper-voice but ALSO by
# jasper-mux/-control (same `jasper` group; mux co-owns StateDirectory=jasper and
# re-chowns the tree to its own user on restart). Without group write they land
# 0644 and the non-owner daemon hits "attempt to write a readonly database" (the
# 2026-06-19 incident). Same directive jasper-fanin/-outputd already use for their
# group-`jasper` sockets. docs/HANDOFF-privilege-separation.md.
UMask=0007
# RuntimeDirectory creates /run/jasper at service start (auto-cleaned
# at stop). voice_daemon binds /run/jasper/voice.sock there for
# jasper-control to send dial hold-to-talk triggers via. With Group=jasper
# the dir is jasper-voice:jasper 0750, so jasper-control (also group jasper)
# can traverse it to reach voice.sock (chmod 0660 in the daemon).
RuntimeDirectory=jasper
RuntimeDirectoryMode=0750

# Hardening — defense in depth around the highest-attack-surface
# process on the speaker (parses microphone audio, third-party LLM
# responses, Spotify/Google API responses). WS1 Phase 3b drops it to the
# non-root jasper-voice user (above); /dev/snd access comes from the `audio`
# supplementary group, not root. WS1 phase-1 hardening (see
# docs/HANDOFF-privilege-separation.md): strict read-only FS (writes pinned by
# ReadWritePaths + StateDirectory/RuntimeDirectory), ProtectHome=tmpfs hides
# /root so a compromised voice can't read root's SSH key to pivot, no kernel
# tunables/modules/logs, no extra namespaces. The capability + syscall filters
# (CapabilityBoundingSet= / SystemCallFilter=@system-service, in the Phase-3b
# block above) carry full value now that the process is non-root.
ProtectSystem=strict
ProtectHome=tmpfs
PrivateTmp=true
NoNewPrivileges=true
# WS1 Phase 4a — jasper-voice only READS /var/lib/jasper-secrets (the LLM keys
# arrive via EnvironmentFile, which systemd reads as root before the drop; the
# Google token tree is read off disk via the jasper-secrets group). It never
# WRITES there — the OAuth save_token / registry.save paths live in jasper-web +
# the operator CLI, not the daemon. ProtectSystem=strict still allows reads, so
# the dir is deliberately NOT in ReadWritePaths: least privilege, and no hard
# dependency on a non-StateDirectory path for the most critical daemon.
# WS1 Phase 4b — Spotify is different: spotipy can refresh and persist OAuth
# tokens from the voice runtime, so the integration-secret compartment is a
# deliberate write path.
ReadWritePaths=/var/lib/jasper /var/lib/jasper-intsecrets
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native

# OOM protection (Stage 1 of memory-resilience plan, see
# docs/HANDOFF-resilience.md). Voice daemon has the LARGEST blast
# radius of any jasper-* process (~150 MB Pss; loads onnx wake model
# + Python deps + holds the LLM session). Conservative protection
# of -500 — kept alive when other processes blow up, but if we
# leak memory ourselves, the kernel should still pick us over
# system-critical infrastructure.
OOMScoreAdjust=-500
# Stage 2 memory bound (the deferred follow-up promised above).
# MemoryHigh (throttle + aggressive reclaim), NOT MemoryMax (kill):
# voice is the most-protected daemon on the box and an outright cgroup
# kill of the recovery-critical wake path is worse than throttling a
# leak while the OOMScoreAdjust ladder + Stage 1 kernel policy decide.
# 384M is ~2.5x the ~140-150 MB Pss steady state in README's "What's
# installed and at what cost" resource table (cgroup memory.current
# counts full RSS + page cache, so headroom over Pss must be generous),
# and matches the per-process MemoryMax budget pi-run-diagnostic.sh
# already treats as safe on a 1 GB Pi. Silent no-op until the memory
# cgroup controller is enabled — same documented-intent posture as
# jasper-mux / -input / -usbsink.
MemoryHigh=384M

[Install]
WantedBy=multi-user.target
