# deploy/systemd/jasper-mux.service
[Unit]
Description=Jasper renderer source-arbiter (auto/manual source selection)
After=jasper-fanin.service librespot.service shairport-sync.service bluealsa.service
Wants=jasper-fanin.service librespot.service shairport-sync.service bluealsa.service

[Service]
Type=simple
EnvironmentFile=-/etc/jasper/jasper.env
EnvironmentFile=-/var/lib/jasper/speaker_name.env
# Spotify credentials are wizard-owned and live outside jasper.env.
# Mux needs them for guarded Spotify handoff: without the Web API router,
# Spotify starts in degraded_safe mode and Camilla stays attenuated.
EnvironmentFile=-/var/lib/jasper/spotify_credentials.env
# Optional fan-in overrides. Loading the same file here keeps
# JASPER_FANIN_CONTROL_SOCKET aligned if an operator moves fan-in's
# UDS out of /run/jasper-fanin/control.sock.
EnvironmentFile=-/var/lib/jasper/fanin.env
ExecStart=/opt/jasper/.venv/bin/jasper-mux
Restart=always
RestartSec=2
# 1 Hz polling of three sources is light; stay under typical
# python+httpx idle footprint. CPU should be <1%.
MemoryHigh=80M
MemoryMax=120M
Nice=-3
RuntimeDirectory=jasper-mux
RuntimeDirectoryMode=0750
# StateDirectory guarantees /var/lib/jasper exists + is writable for the
# durable source-selection mode file (mux_mode.json), so a household's
# manual source pin survives the Restart=always deploy/restart cycle.
# Mux runs as root with no ProtectSystem today, so /var/lib/jasper is
# already writable; this codifies the write target and keeps the pin
# durable if mux ever gains filesystem sandboxing. Mirrors jasper-voice.
StateDirectory=jasper
StateDirectoryMode=0750

# OOM protection (Stage 1 of memory-resilience plan, see
# docs/HANDOFF-resilience.md). Mux is restartable but now sits on the
# audible source-selection path: fan-in starts safe/closed and mux
# opens the chosen renderer lane after volume handoff. Keep moderate
# protection, but restart on every exit.
# NOTE: MemoryHigh/MemoryMax above are silent no-ops until the
# memory cgroup is enabled (Stage 2 of plan). Kept here as
# documentation of intent.
OOMScoreAdjust=-300

[Install]
WantedBy=multi-user.target
