# deploy/systemd/jasper-usbgadget.service
[Unit]
Description=Jasper composite USB gadget (management network + optional audio)
# The single owner of the ConfigFS composite gadget. It composes up to two
# functions onto the one dwc2 UDC:
#   - ncm.usb0  — default-on USB management network where the resolved
#     hardware transport is available (a laptop plugged into the USB-C port
#     reaches http://<hostname>/ even with Wi-Fi off);
#   - uac2.usb0 — the wizard-toggled USB audio input (jasper-usbsink source).
#
# This is the FIRST gadget unit install.sh enables (systemctl enable --now),
# deliberately — it carries the management network. jasper-usbsink.service is a
# process-free readiness marker that orders After= this unit and does not own the
# gadget.
#
# There is intentionally NO jasper-local-source-allowed ExecCondition on this
# unit: the network function must survive multiroom follower parking. The
# AUDIO-function gate lives inside jasper-usbgadget-up: it reads canonical USB
# source intent plus role permission directly and composes uac2 only when both
# allow it. jasper-usbsink enablement remains a derived lifecycle mirror, and
# the service keeps the same source-aware ExecCondition at its start boundary.
#
# NAMING DEBT: the ConfigFS gadget dir name stays `jts-usb-audio` (the doctor's
# low-latency contract reads that exact path). Accepted debt — see
# docs/HANDOFF-usb-gadget.md.
After=systemd-modules-load.service sys-kernel-config.mount jasper-audio-hardware-reconcile.service jasper-fanin.service
Wants=jasper-audio-hardware-reconcile.service jasper-fanin.service
Before=jasper-usbsink.service
ConditionPathExists=/sys/kernel/config

[Service]
Type=oneshot
RemainAfterExit=yes
# ConfigFS composition is synchronous and normally sub-second. Bound both
# halves so a coordinator `restart` has a 10 s service contract; its client
# waits one second longer for PID 1's terminal result.
TimeoutStartSec=5s
TimeoutStopSec=5s
EnvironmentFile=-/etc/jasper/jasper.env
# The management-writable speaker-name file is parsed as one inert value
# by gadget-up; never inject it into this root process. Test-only path/command
# seams are also stripped from the production unit environment.
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
UnsetEnvironment=LD_PRELOAD LD_LIBRARY_PATH LD_AUDIT GLIBC_TUNABLES PYTHONPATH PYTHONHOME JASPER_CONFIGFS_ROOT JASPER_PI_MODEL_FILE JTS_BOOT_CONFIG_FILE JASPER_UDC_CLASS_DIR JASPER_CPUINFO_FILE JASPER_MODULES_ROOT JASPER_SPEAKER_NAME_FILE JASPER_SPEAKER_NAME_READER JASPER_USBGADGET_AUDIO_ALLOWED_CMD JASPER_USBGADGET_AUDIO_READY_CMD JASPER_USBGADGET_AUDIO_DATA_READY_CMD JASPER_USBGADGET_HARDWARE_ALLOWED_CMD
# Skip (NOT fail) when hardware transport is unavailable, there is no UDC yet,
# or no function is wanted (network kill-switched AND audio off/parked). The
# zero-function case restores the historic zero-RAM contract: libcomposite
# never loads. Exit 1 from an ExecCondition is systemd's "condition not met".
ExecCondition=/usr/local/sbin/jasper-usbgadget-wanted
# modprobe libcomposite drags in the gadget framework + u_ether/u_audio deps.
# A failure means the kernel lacks CONFIG_USB_GADGET / CONFIG_USB_CONFIGFS —
# surface that as a unit failure rather than a mysterious later error.
ExecStartPre=/sbin/modprobe libcomposite
# Make the host-visible audio device name track the Speaker Name by patching
# the kernel's hardcoded UAC2 AudioStreaming string. Best-effort (leading `-`):
# a failure only means the host shows the default label — it must never block
# the gadget. Runs before gadget-up so the patched module autoloads when the
# uac2 function is created. See docs/HANDOFF-usbsink.md "Device name".
ExecStartPre=-/usr/local/sbin/jasper-usbsink-name-patch
ExecStart=/usr/local/sbin/jasper-usbgadget-up
ExecStop=/usr/local/sbin/jasper-usbgadget-down
# Best-effort rmmod in dependency order: the network function modules
# (usb_f_ncm, u_ether) and the audio function modules (usb_f_uac2, u_audio)
# ahead of libcomposite. Leading `-` ignores non-zero exit — if a module is
# still in use (stale descriptor, stuck capture) the rmmod fails harmlessly.
# Either way the gadget descriptor is gone, so the host stops seeing JTS.
ExecStopPost=-/sbin/rmmod usb_f_ncm
ExecStopPost=-/sbin/rmmod usb_f_uac2
ExecStopPost=-/sbin/rmmod u_ether
ExecStopPost=-/sbin/rmmod u_audio
ExecStopPost=-/sbin/rmmod libcomposite
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
