# deploy/systemd/jasper-usbsink-volume.service
[Unit]
Description=Jasper USB sink host-volume observer
# Non-real-time helper for jasper-usbsink.service. The Rust bridge owns the
# audio data plane and preempt endpoint; this Python process only polls the
# UAC2 gadget mixer and reports host slider changes to jasper-control.
#
# It tracks the AUDIO daemon's lifecycle, not the gadget's: the UAC2 mixer it
# polls only exists while the uac2 function is composed (i.e. while
# jasper-usbsink is running). The composite gadget now OUTLIVES audio (it also
# carries the always-on network), so ordering/binding to jasper-usbsink.service
# — which itself Requires= the gadget — keeps this helper off whenever audio
# is parked while the network stays up.
After=jasper-usbsink.service jasper-control.service
Requires=jasper-usbsink.service
PartOf=jasper-usbsink.service

[Service]
Type=simple
EnvironmentFile=-/etc/jasper/jasper.env
EnvironmentFile=-/var/lib/jasper/usbsink.env
ExecCondition=/opt/jasper/.venv/bin/jasper-local-source-allowed
ExecStart=/opt/jasper/.venv/bin/jasper-usbsink-volume
Restart=on-failure
RestartSec=2s
Nice=5
OOMScoreAdjust=100
MemoryHigh=32M
MemoryMax=64M

[Install]
WantedBy=multi-user.target
