# deploy/systemd/jasper-wifi-recover.service
[Unit]
Description=Recover Wi-Fi after transient AP or Pi 5 scan-suppression loss
Documentation=file:/opt/jasper/docs/HANDOFF-resilience.md
After=NetworkManager.service
Wants=NetworkManager.service
ConditionPathExists=/var/lib/jasper/wifi_guardian.env

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/jasper-wifi-recover --reason systemd
# Bounded: active Wi-Fi exits immediately; fault path may run one scan-repair
# attempt plus one guardian activation/recreate attempt.
TimeoutStartSec=75
StandardOutput=journal
StandardError=journal
