# deploy/systemd/librespot.service
[Unit]
Description=Spotify Connect renderer (librespot, log volume curve)
Documentation=https://github.com/librespot-org/librespot
After=network-online.target sound.target jasper-camilla.service
Wants=network-online.target
# Always-on renderer: same logic as shairport-sync — Spotify Connect
# must be advertised on the LAN whenever the speaker is up, including
# after a clean exit (which Restart=on-failure ignores by design).
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=pi
Group=audio
# Spotify Connect display name. Sourced from /etc/jasper/jasper.env if
# JASPER_SPOTIFY_DEVICE_NAME is set there (the EnvironmentFile= below);
# otherwise falls back to the Environment= default. This keeps a second
# Pi on the same LAN from crashing librespot with an Avahi name
# collision — install.sh seeds jasper.env with a hostname-derived
# default via jasper-derive-device-name (jts2 → "JTS-2", etc.). The
# leading `-` on EnvironmentFile makes it optional (no error if the
# file is absent during very early bringup).
Environment="JASPER_SPOTIFY_DEVICE_NAME=JTS"
EnvironmentFile=-/etc/jasper/jasper.env
# RuntimeDirectory creates /run/librespot/ at start with mode 0755.
# Other components (jasper-mux, jasper-voice volume_observers,
# jasper-control coordinator) read /run/librespot/state.json — the
# file is written atomically by the --onevent hook on every player
# event (volume_set, playing, paused, stopped, etc).
RuntimeDirectory=librespot
RuntimeDirectoryMode=0755
# CacheDirectory creates /var/cache/librespot/ owned by pi:audio.
# Used by --system-cache below for credential persistence — without
# this, librespot runs in pure zeroconf mode and is invisible to the
# Spotify Web API's `devices()` list until a phone manually claims it
# via Spotify Connect, breaking voice-driven `spotify_play` cold-starts
# any time librespot has restarted since the last claim.
CacheDirectory=librespot
CacheDirectoryMode=0700
# Volume curve: --volume-ctrl log + --volume-range 60 gives a
# perceptually linear slider over 60 dB. go-librespot's hardcoded
# cubic curve packed all useful loudness into the bottom third —
# unusable on real speakers. See docs/HANDOFF-volume.md.
#
# Format: S24_3 + tpdf dither = 24-bit packed output with triangular-
# distribution dither, ~141 dB SNR, well below any DAC noise floor.
#
# Audio device: jasper_renderer_in — the plug-wrapped renderer-side
# dmix defined in /root/.asoundrc (deploy/alsa/asoundrc.jasper).
# All three renderers (librespot, shairport-sync, bluealsa-aplay)
# write to this same logical device; dmix multi-writer-safely sums
# their streams onto snd-aloop. Before 2026-05-22 the renderers
# wrote directly to plughw:Loopback,0,0 (single-writer), which
# caused librespot to crash-loop on ALSA -EBUSY when shairport had
# the device open during phantom AirPlay SETUPs. See the rationale
# block at the top of the asoundrc file.
ExecStart=/usr/bin/librespot \
    --name ${JASPER_SPOTIFY_DEVICE_NAME} \
    --device-type speaker \
    --bitrate 320 \
    --backend alsa \
    --device jasper_renderer_in \
    --format S24_3 \
    --dither tpdf \
    --mixer softvol \
    --volume-ctrl log \
    --volume-range 60 \
    --initial-volume 50 \
    --enable-volume-normalisation \
    --system-cache /var/cache/librespot \
    --disable-audio-cache \
    --onevent /usr/local/bin/jasper-librespot-event
Restart=always
RestartSec=2
MemoryHigh=150M
MemoryMax=200M
Nice=-5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
