# deploy/systemd/librespot.service
[Unit]
Description=Spotify Connect renderer (librespot, log volume curve)
Documentation=https://github.com/librespot-org/librespot
After=network-online.target sound.target jasper-fanin.service jasper-camilla.service
Wants=network-online.target jasper-fanin.service
# Always-on renderer: same logic as shairport-sync — Spotify Connect
# must be advertised on the LAN whenever the speaker is up, including
# after a clean exit (which Restart=on-failure ignores by design).
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
# Stage 2 audio-protection: pages never swap to zram. See
# docs/HANDOFF-resilience.md "Stage 2 — audio-protection subset".
Slice=jts-audio.slice
Type=simple
User=pi
Group=audio
# Spotify Connect display name. Sourced from /etc/jasper/jasper.env if
# JASPER_SPEAKER_NAME is the one user-facing display name shown by
# Spotify Connect, AirPlay, Bluetooth, and USB Audio. The /speaker/
# wizard writes /var/lib/jasper/speaker_name.env; the Environment=
# default covers very early bringup before that file exists.
Environment="JASPER_SPEAKER_NAME=JTS"
EnvironmentFile=-/etc/jasper/jasper.env
EnvironmentFile=-/var/lib/jasper/speaker_name.env
# RuntimeDirectory creates /run/librespot/ at start with mode 0755.
# Other components (jasper-mux, jasper-voice volume_observers,
# jasper-control coordinator) read /run/librespot/state.json — the
# file is written atomically by the --onevent hook on every player
# event (volume_set, playing, paused, stopped, etc).
RuntimeDirectory=librespot
RuntimeDirectoryMode=0755
# CacheDirectory creates /var/cache/librespot/ owned by pi:audio.
# Used by --system-cache below for credential persistence — without
# this, librespot runs in pure zeroconf mode and is invisible to the
# Spotify Web API's `devices()` list until a phone manually claims it
# via Spotify Connect, breaking voice-driven `spotify_play` cold-starts
# any time librespot has restarted since the last claim.
CacheDirectory=librespot
CacheDirectoryMode=0700
# Volume curve: --volume-ctrl log + --volume-range 60 gives a
# perceptually linear slider over 60 dB. go-librespot's hardcoded
# cubic curve packed all useful loudness into the bottom third —
# unusable on real speakers. See docs/HANDOFF-volume.md.
#
# Format: S24_3 + tpdf dither = 24-bit packed output with triangular-
# distribution dither, ~141 dB SNR, well below any DAC noise floor.
#
# Audio device: librespot_substream — Spotify's private snd-aloop
# lane. jasper-fanin reads the capture side and publishes the summed
# music stream to CamillaDSP/AEC. The private lane keeps renderer
# opens independent; latest-source-wins remains jasper-mux policy.
ExecStart=/usr/bin/librespot \
    --name=${JASPER_SPEAKER_NAME} \
    --device-type speaker \
    --bitrate 320 \
    --backend alsa \
    --device librespot_substream \
    --format S24_3 \
    --dither tpdf \
    --mixer softvol \
    --volume-ctrl log \
    --volume-range 60 \
    --initial-volume 50 \
    --enable-volume-normalisation \
    --system-cache /var/cache/librespot \
    --disable-audio-cache \
    --onevent /usr/local/bin/jasper-librespot-event
Restart=always
RestartSec=2
MemoryHigh=150M
MemoryMax=200M
Nice=-5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
