# deploy/systemd/nqptp.service
# Vendored from mikebrady/nqptp's `make install` output. Identical to
# what nqptp's autotools build installs at
# /usr/local/lib/systemd/system/nqptp.service — captured here so
# install.sh can deploy it without depending on whether the source
# build's install-systemd target succeeded (it sometimes fails on
# Trixie due to an `install` flag mismatch — the binary still installs
# fine to /usr/local/bin/nqptp but the unit step errors out).

[Unit]
Description=NQPTP -- Not Quite PTP
After=network.target
# nqptp is a hard Requires= of shairport-sync (AirPlay 2 timing). While AirPlay
# is household-On and role-authorized, keep the timing companion running across
# crashes and clean exits. Its matching ExecCondition prevents Restart=always
# from overriding household Off or follower parking.
StartLimitIntervalSec=600
StartLimitBurst=20

[Service]
# Type=simple has no readiness wait; finite ceilings keep coordinator jobs
# bounded and are paired with slightly longer client waits.
TimeoutStartSec=2s
TimeoutStopSec=5s
UnsetEnvironment=LD_PRELOAD LD_LIBRARY_PATH LD_AUDIT GLIBC_TUNABLES
ExecCondition=+/usr/bin/env -i PATH=/opt/jasper/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin /opt/jasper/.venv/bin/jasper-local-source-allowed --source airplay
ExecStart=/usr/local/bin/nqptp
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
