# deploy/systemd/nqptp.service
# Vendored from mikebrady/nqptp's `make install` output. Identical to
# what nqptp's autotools build installs at
# /usr/local/lib/systemd/system/nqptp.service — captured here so
# install.sh can deploy it without depending on whether the source
# build's install-systemd target succeeded (it sometimes fails on
# Trixie due to an `install` flag mismatch — the binary still installs
# fine to /usr/local/bin/nqptp but the unit step errors out).

[Unit]
Description=NQPTP -- Not Quite PTP
After=network.target
# nqptp is a hard Requires= of shairport-sync (AirPlay 2 timing). If
# it dies cleanly, shairport stays up but every new session fails to
# negotiate the PTP clock. Match shairport's always-on stance: keep
# nqptp running across crashes AND clean exits.
StartLimitIntervalSec=600
StartLimitBurst=20

[Service]
ExecCondition=/opt/jasper/.venv/bin/jasper-local-source-allowed
ExecStart=/usr/local/bin/nqptp
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
