# deploy/systemd/shairport-sync.service
[Unit]
Description=Shairport Sync AirPlay 2 audio receiver (JTS)
After=sound.target jasper-fanin.service avahi-daemon.service nqptp.service network-online.target
Requires=avahi-daemon.service nqptp.service
Wants=network-online.target jasper-fanin.service
# Always-on renderer: the AirPlay endpoint MUST be advertised on the
# LAN whenever the speaker is up. A clean exit (status=0) — whether
# from a SIGTERM during a shell session, an admin `systemctl stop`
# that wasn't followed by a `start`, or a clean shutdown path inside
# shairport itself — would, with Restart=on-failure, leave the device
# silently invisible to AirPlay clients until the next reboot.
# Restart=always covers crash AND clean-exit paths. Start-limit guard
# stops genuine config breakage from infinite-restart-looping.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
# Stage 2 audio-protection: pages never swap to zram. See
# docs/HANDOFF-resilience.md "Stage 2 — audio-protection subset".
Slice=jts-audio.slice
Type=simple
# Render /etc/shairport-sync.conf from the template every time the
# service starts, picking up the current value of
# /var/lib/jasper/airplay_mode.env. The leading `+` lets the
# pre-start helper run as root so it can write /etc; ExecStart drops
# back to the unprivileged shairport-sync user. The renderer is
# idempotent and fast (sed substitution), so running it on every
# start has negligible cost.
ExecStartPre=+/usr/local/sbin/jasper-apply-airplay-mode
ExecStart=/usr/local/bin/shairport-sync
User=shairport-sync
Group=shairport-sync
SupplementaryGroups=audio
# Scheduling priority matches jasper-camilla. shairport produces audio
# that camilla consumes via snd-aloop; if shairport runs at default
# Nice=0 while camilla + jasper-aec-bridge run at Nice=-10, brief
# scheduler stalls on shairport let camilla drain the loopback and
# underrun. Matching priorities removes the asymmetry.
Nice=-10
IOSchedulingClass=realtime
IOSchedulingPriority=4
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
