# deploy/systemd/bluetooth.service.d/jts-timeout.conf
# bluetooth.service is shared control-plane infrastructure, but the source
# coordinator may start it synchronously before enabling Bluetooth audio.
# Bound the package unit so the coordinator's client cannot return while PID 1
# may still legally be running the same job.
[Service]
TimeoutStartSec=10s
TimeoutStopSec=10s
