# deploy/jasper-sources-web.service
[Unit]
Description=Jasper source toggles (web at /sources/)
# Endpoint-safe standalone host for /sources/. Full speakers serve the same
# page from the combined jasper-web process; endpoint installs use this unit
# so they do not import the full setup-wizard bundle.
Requires=jasper-sources-web.socket
After=jasper-sources-web.socket network-online.target
Wants=network-online.target

[Service]
Type=notify
WatchdogSec=30s
EnvironmentFile=-/etc/jasper/jasper.env
ExecStart=/opt/jasper/.venv/bin/jasper-sources-web --host 127.0.0.1 --port 8773
Restart=on-failure
RestartSec=2

# Hardening. Runs as root because toggles call systemctl/BlueZ DBus for
# operator-approved source state changes. The service binds only through
# the loopback socket owned by systemd and nginx fronts it on jts.local.
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

# Small stdlib http.server plus occasional systemctl probes.
MemoryHigh=60M
MemoryMax=90M

[Install]
# Enable the .socket, not the .service.
