# deploy/jasper-web.socket
[Unit]
Description=Jasper web setup wizards (Spotify, voice, Google, AirPlay, Sources, Speaker, Wake, Wi-Fi, Transit, Home Assistant, Weather, Sound, Wake-Corpus, Rooms) — socket
# When any of these ports gets a connection, systemd will activate
# jasper-web.service to handle it. The service exits after 10 min
# of no activity; this .socket stays listening throughout, so
# nginx-fronted requests transparently re-spawn the daemon.

[Socket]
# 14 ListenStream= → 14 fds passed to the daemon at LISTEN_FDS_START.
# Order is undefined per systemd issue #16647, so the Python side
# matches by getsockname()[1] (port) rather than by fd index.
# All ports bind 127.0.0.1 — multi-Pi households need no cross-host
# port coordination since each Pi has its own loopback.
ListenStream=127.0.0.1:8765
ListenStream=127.0.0.1:8767
ListenStream=127.0.0.1:8768
ListenStream=127.0.0.1:8771
ListenStream=127.0.0.1:8773
ListenStream=127.0.0.1:8774
ListenStream=127.0.0.1:8775
ListenStream=127.0.0.1:8777
ListenStream=127.0.0.1:8778
ListenStream=127.0.0.1:8779
# 8782 = wake-corpus recorder (jasper.web.wake_corpus_setup). Longer-
# lived than the other wizards — recording sessions can run an hour
# while the operator works through 60+ utterances. Browser polls
# /api/status every 2s, so the idle tracker stays bumped while a
# tab is open.
ListenStream=127.0.0.1:8782
# 8783 = speaker display-name wizard (jasper.web.speaker_setup).
ListenStream=127.0.0.1:8783
# 8784 = sound curve + preference EQ wizard (jasper.web.sound_setup).
ListenStream=127.0.0.1:8784
# 8785 = room/group directory wizard (jasper.web.rooms_setup).
ListenStream=127.0.0.1:8785
# Accept=no — single daemon for all fds; Python startup cost forbids
# per-connection accept.
Accept=no
# Default backlog (128) is fine for a setup wizard.
# Bump TriggerLimitBurst above the systemd default (20) — dev reload
# loops + nginx healthchecks can otherwise trip the rate-limit and
# leave the socket in failure state requiring `reset-failed`.
TriggerLimitIntervalSec=10s
TriggerLimitBurst=100

[Install]
WantedBy=sockets.target
