# deploy/jasper-web.service
[Unit]
Description=Jasper smart-speaker web setup service (Spotify + voice provider + Google + AirPlay + Sources + Speaker + Wake + Wi-Fi + Transit + Weather + Sound + Rooms + Tools)
# Socket-activated by jasper-web.socket — daemon exits after 10 min idle.
Requires=jasper-web.socket
After=jasper-web.socket network-online.target time-sync.target
Wants=network-online.target time-sync.target
# Socket-activated but recoverable: if a request arrives during transient
# memory pressure and the daemon is OOM-killed repeatedly, do not park the
# wizard surface after systemd's tiny default burst.
StartLimitIntervalSec=600
StartLimitBurst=20
# time-sync.target gate: Pi 5 has no battery-backed RTC, so on first
# boot the clock can be in 1970 until NTP catches up. The Spotify
# (and Google) token-exchange POSTs fail cert validation when the
# system clock is too far off, and the user-visible symptom would be
# "Auth exchange failed: certificate has expired" right at the end
# of an otherwise-successful OAuth flow. Waiting on time-sync.target
# makes that race impossible.

[Service]
# Type=notify — daemon sends READY=1 after binding so nginx's
# proxy_pass doesn't see a half-up service. WATCHDOG=1 every ~15s
# keeps systemd happy; a hang gets killed and re-spawned by the
# next incoming request.
Type=notify
WatchdogSec=30s
EnvironmentFile=/etc/jasper/jasper.env
EnvironmentFile=-/var/lib/jasper/speaker_name.env
# Optional file written by the Spotify setup wizard when the user
# pastes their Spotify Developer App credentials. Sourced after
# jasper.env so wizard-written values win. Leading `-` makes it
# optional — fresh installs boot fine without it (and land on the
# wizard).
# WS1 Phase 4b — moved into the group-`jasper-intsecrets` dir alongside
# the Spotify OAuth token cache.
EnvironmentFile=-/var/lib/jasper-intsecrets/spotify_credentials.env
# Optional file written by the voice-provider wizard at
# /voice/. Same override pattern: wizard values beat operator
# defaults. Sourced here so the wizard's process can render the
# current state correctly even on cold boot. WS1 Phase 4a: this holds
# only the non-secret selectors now; the API keys are in voice_keys.env.
EnvironmentFile=-/var/lib/jasper/voice_provider.env
# WS1 Phase 4a — the GEMINI/OPENAI/XAI API keys (group `jasper-secrets`).
# The /voice/ wizard writes + renders them; jasper-web joins the
# jasper-secrets group (below) so it can read/replace this file. systemd
# itself reads EnvironmentFile= as root, so sourcing it always works.
EnvironmentFile=-/var/lib/jasper-secrets/voice_keys.env
# Optional file written by the Google wizard at /google/ — holds
# GOOGLE_CLIENT_ID / GOOGLE_CLIENT_SECRET. Same override pattern.
# WS1 Phase 4a — moved into the group-`jasper-secrets` dir.
EnvironmentFile=-/var/lib/jasper-secrets/google_credentials.env
# Google Routes API key written by /transit/. This is a billable Maps API
# key, so it lives in the jasper-secrets compartment rather than transit.env.
EnvironmentFile=-/var/lib/jasper-secrets/google_routes.env
# Optional file written by the wake-word picker at /wake/ — holds
# JASPER_WAKE_MODEL pointing at either a bundled openWakeWord name
# or an absolute path under /var/lib/jasper/wake/. The wizard reads
# this file to render which model is active and writes back here on
# save. Sourced AFTER jasper.env so wizard values win.
EnvironmentFile=-/var/lib/jasper/wake_model.env
# Optional file written by the weather wizard at /weather/ — holds
# rounded coordinates + temperature units for default-location weather
# questions. Sourced AFTER jasper.env so wizard values win.
EnvironmentFile=-/var/lib/jasper/weather.env
# Optional file written by the transit wizard at /transit/ — holds
# geocoded coordinates + subway/bus stop selections + the BusTime API
# key. Same override pattern: sourced AFTER jasper.env so the wizard's
# picks win over operator-managed defaults.
EnvironmentFile=-/var/lib/jasper/transit.env
# Single Python process serves all wizards as separate threads:
#   /spotify/  →  127.0.0.1:8765
#   /voice/    →  127.0.0.1:8767
#   /google/   →  127.0.0.1:8768
#   /airplay/  →  127.0.0.1:8771
#   /sources/  →  127.0.0.1:8773
#   /speaker/  →  127.0.0.1:8783
#   /wake/     →  127.0.0.1:8774
#   /wifi/     →  127.0.0.1:8775
#   /transit/  →  127.0.0.1:8777
#   /ha/       →  127.0.0.1:8778
#   /weather/  →  127.0.0.1:8779
#   /sound/    →  127.0.0.1:8784
#   /rooms/    →  127.0.0.1:8785
# When socket-activated, all fds arrive via LISTEN_FDS and the
# Python side disambiguates by getsockname()[1] (port). See
# jasper/web/__main__.py for the wiring.
ExecStart=/opt/jasper/.venv/bin/python -m jasper.web
# Restart=on-failure (NOT always) so clean idle-exit stays exited —
# the .socket re-spawns us on next request. Restart=always would
# defeat the whole memory-saving point.
Restart=on-failure
RestartSec=2

# WS1 Phase 3b — user drop (3b-3). jasper-web is the most network-exposed
# surface (untrusted SSID / form / OAuth input), so de-rooting it is the highest-
# value drop — and the riskiest, because a Wi-Fi-management break is the
# worst-case brick for a headless, often Ethernet-less speaker. It runs as a
# dedicated non-root service user in the shared `jasper` group (the cross-daemon
# /run socket + /var/lib/jasper access boundary).
#
# This drop is GATED on hardware validation (recovery matrix in
# docs/HANDOFF-privilege-separation.md) of every privileged surface it owns.
# Each one's mechanism, validated on jts.local (systemd 257, polkit 126, NM 1.52):
#   * NetworkManager (/wifi/ scan/connect/forget/radio/PSK-stash) — a polkit
#     rule (deploy/polkit/49-jasper-web.rules) grants this uid the 5 NM actions
#     wifi_setup.py drives; NM denies a sessionless daemon by default, so the
#     rule is load-bearing. `netdev` group is neither necessary nor sufficient.
#   * BlueZ adapter Alias/Powered (/speaker rename, BT radio) — the `bluetooth`
#     supplementary group (D-Bus policy grant, not polkit).
#   * journalctl -k (Wi-Fi scan-suppression diagnostics) — the `systemd-journal`
#     supplementary group (also fail-soft, so non-load-bearing).
#   * NL80211 scan-repair (jasper/wifi_scan_repair.py) — intentionally NOT given
#     CAP_NET_ADMIN; /wifi/scan asks jasper-control's broker to start the
#     root-only jasper-wifi-scan-repair.service helper instead.
#   * systemctl restarts after a wizard save — already routed through
#     jasper-control's restart broker (manage_units); once non-root the broker
#     is the only path (geteuid()!=0 removes the direct-systemctl fallback).
#
# CapabilityBoundingSet= drops ALL caps (the privileged ops go through polkit /
# the broker / D-Bus groups, never a capability — and the scan-repair cap is
# deliberately withheld from this web process). SystemCallFilter=@system-service
# blocks the dangerous syscall classes while permitting fork/execve, so the nmcli / busctl /
# journalctl children it spawns exec cleanly. ProtectKernelLogs stays OMITTED
# (below) for the same diagnostic-subprocess caution as jasper-control.
User=jasper-web
Group=jasper
# WS1 Phase 4a — `jasper-secrets` (supplementary) lets the wizards write
# + read the split-out LLM keys, Google client secret, and Google OAuth
# token tree in /var/lib/jasper-secrets. WS1 Phase 4b — `jasper-intsecrets`
# lets the wizards write/read HA + Spotify credentials and token caches.
# `audio` (active-speaker commissioning correction_substream tone),
# `bluetooth` (BlueZ Alias), `systemd-journal` (journalctl -k).
# docs/HANDOFF-privilege-separation.md.
SupplementaryGroups=audio bluetooth systemd-journal jasper-secrets jasper-intsecrets
CapabilityBoundingSet=
SystemCallFilter=@system-service

# Hardening — this process needs to read /etc/jasper/jasper.env,
# read+write /var/lib/jasper/ and generated CamillaDSP configs,
# update BlueZ's local name (/etc/bluetooth/main.conf), re-render the speaker-
# name Avahi advert (/etc/avahi/services), accept connections on systemd-passed
# loopback fds, and drive the restart broker after a wizard updates state.
# /etc/bluetooth and /var/lib/camilladsp/configs are root:jasper 2775 (setgid,
# group-writable, set by install.sh) so the non-root user can atomically replace
# files there; /etc/avahi/services is the same (set by 3b-2's install). Adding
# /etc/avahi/services to ReadWritePaths also fixes a latent bug: under
# ProtectSystem=strict the /speaker rename's Avahi re-render silently no-op'd
# (the dir was outside the writable set).
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
# /var/lib/jasper-secrets (WS1 Phase 4a) is writable so the /voice +
# /google wizards can atomically replace the secret files + write OAuth
# tokens there (setgid dir → group `jasper-secrets`).
# /var/lib/jasper-intsecrets (WS1 Phase 4b) is writable so the /spotify +
# /ha wizards can atomically replace integration secret files and Spotify
# OAuth tokens there.
ReadWritePaths=/var/lib/jasper /var/lib/jasper-secrets /var/lib/jasper-intsecrets /var/lib/camilladsp/configs /etc/bluetooth /etc/avahi/services
# WS1 phase-1 hardening (docs/HANDOFF-privilege-separation.md): jasper-web gets
# the full phase-1 stanza. ProtectKernelLogs is OMITTED — the wizards shell to
# nmcli / journalctl -k for Wi-Fi diagnostics and must not be starved of kernel-
# log access. CapabilityBoundingSet= + SystemCallFilter= landed with the 3b-3
# user drop (above).
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
SystemCallArchitectures=native

[Install]
# WantedBy left out intentionally — enabling the .socket pulls in
# the service on demand. `systemctl enable jasper-web.socket`.
