# deploy/jasper-dial-web.socket
[Unit]
Description=Jasper dial onboarding wizard — socket activation

[Socket]
ListenStream=127.0.0.1:8766
Accept=no
TriggerLimitIntervalSec=10s
TriggerLimitBurst=100

[Install]
WantedBy=sockets.target
