# deploy/jasper-system-web.socket
[Unit]
Description=Jasper system dashboard — socket activation

[Socket]
ListenStream=127.0.0.1:8772
Accept=no
TriggerLimitIntervalSec=10s
TriggerLimitBurst=100

[Install]
WantedBy=sockets.target
