# deploy/jasper-web.service
[Unit]
Description=Jasper smart-speaker web setup service (Spotify + voice provider + Google + AirPlay + Sources + Speaker + Wake + Wi-Fi + Transit + Weather + Sound + Rooms)
# Socket-activated by jasper-web.socket — daemon exits after 10 min idle.
Requires=jasper-web.socket
After=jasper-web.socket network-online.target time-sync.target
Wants=network-online.target time-sync.target
# time-sync.target gate: Pi 5 has no battery-backed RTC, so on first
# boot the clock can be in 1970 until NTP catches up. The Spotify
# (and Google) token-exchange POSTs fail cert validation when the
# system clock is too far off, and the user-visible symptom would be
# "Auth exchange failed: certificate has expired" right at the end
# of an otherwise-successful OAuth flow. Waiting on time-sync.target
# makes that race impossible.

[Service]
# Type=notify — daemon sends READY=1 after binding so nginx's
# proxy_pass doesn't see a half-up service. WATCHDOG=1 every ~15s
# keeps systemd happy; a hang gets killed and re-spawned by the
# next incoming request.
Type=notify
WatchdogSec=30s
EnvironmentFile=/etc/jasper/jasper.env
EnvironmentFile=-/var/lib/jasper/speaker_name.env
# Optional file written by the Spotify setup wizard when the user
# pastes their Spotify Developer App credentials. Sourced after
# jasper.env so wizard-written values win. Leading `-` makes it
# optional — fresh installs boot fine without it (and land on the
# wizard).
EnvironmentFile=-/var/lib/jasper/spotify_credentials.env
# Optional file written by the voice-provider wizard at
# /voice/. Same override pattern: wizard values beat operator
# defaults. Sourced here so the wizard's process can render the
# current state correctly even on cold boot.
EnvironmentFile=-/var/lib/jasper/voice_provider.env
# Optional file written by the Google wizard at /google/ — holds
# GOOGLE_CLIENT_ID / GOOGLE_CLIENT_SECRET. Same override pattern.
EnvironmentFile=-/var/lib/jasper/google_credentials.env
# Optional file written by the wake-word picker at /wake/ — holds
# JASPER_WAKE_MODEL pointing at either a bundled openWakeWord name
# or an absolute path under /var/lib/jasper/wake/. The wizard reads
# this file to render which model is active and writes back here on
# save. Sourced AFTER jasper.env so wizard values win.
EnvironmentFile=-/var/lib/jasper/wake_model.env
# Optional file written by the weather wizard at /weather/ — holds
# rounded coordinates + temperature units for default-location weather
# questions. Sourced AFTER jasper.env so wizard values win.
EnvironmentFile=-/var/lib/jasper/weather.env
# Optional file written by the transit wizard at /transit/ — holds
# geocoded coordinates + subway/bus stop selections + the BusTime API
# key. Same override pattern: sourced AFTER jasper.env so the wizard's
# picks win over operator-managed defaults.
EnvironmentFile=-/var/lib/jasper/transit.env
# Single Python process serves all wizards as separate threads:
#   /spotify/  →  127.0.0.1:8765
#   /voice/    →  127.0.0.1:8767
#   /google/   →  127.0.0.1:8768
#   /airplay/  →  127.0.0.1:8771
#   /sources/  →  127.0.0.1:8773
#   /speaker/  →  127.0.0.1:8783
#   /wake/     →  127.0.0.1:8774
#   /wifi/     →  127.0.0.1:8775
#   /transit/  →  127.0.0.1:8777
#   /ha/       →  127.0.0.1:8778
#   /weather/  →  127.0.0.1:8779
#   /sound/    →  127.0.0.1:8784
#   /rooms/    →  127.0.0.1:8785
# When socket-activated, all fds arrive via LISTEN_FDS and the
# Python side disambiguates by getsockname()[1] (port). See
# jasper/web/__main__.py for the wiring.
ExecStart=/opt/jasper/.venv/bin/python -m jasper.web
# Restart=on-failure (NOT always) so clean idle-exit stays exited —
# the .socket re-spawns us on next request. Restart=always would
# defeat the whole memory-saving point.
Restart=on-failure
RestartSec=2

# Hardening — this process needs to read /etc/jasper/jasper.env,
# read+write /var/lib/jasper/ and generated CamillaDSP configs,
# update BlueZ's local name, accept connections on systemd-passed
# loopback fds, and call `systemctl restart` after a wizard updates
# state.
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/var/lib/jasper /var/lib/camilladsp/configs /etc/bluetooth

[Install]
# WantedBy left out intentionally — enabling the .socket pulls in
# the service on demand. `systemctl enable jasper-web.socket`.
