# deploy/systemd/bluealsa.service.d/jts-restart.conf
# Authorized-renderer hardening: while household Bluetooth intent is On and the
# role permits local sources, bluealsa hosts the BlueZ A2DP backend. Its package
# default Restart= policy
# is on-failure (or whatever the package default is), which silently
# drops Bluetooth audio after a clean exit (status=0) — same failure
# mode that PR #4 fixed for shairport-sync, librespot, nqptp, and
# jasper-camilla. Override unconditionally so a clean exit auto-
# recovers; the source-aware ExecCondition below prevents household Off or role
# parking from being overridden. Keep the start-limit window long enough to ride out transient
# memory pressure, but finite so genuine config breakage still parks loudly.
[Unit]
StartLimitIntervalSec=600
StartLimitBurst=20

[Service]
# bluealsa is D-Bus-backed, so allow a short registration window while keeping
# the source coordinator's blocking job finite.
TimeoutStartSec=5s
TimeoutStopSec=5s
UnsetEnvironment=LD_PRELOAD LD_LIBRARY_PATH LD_AUDIT GLIBC_TUNABLES
ExecCondition=+/usr/bin/env -i PATH=/opt/jasper/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin /opt/jasper/.venv/bin/jasper-local-source-allowed --source bluetooth
Restart=always
RestartSec=5
