# deploy/systemd/bluealsa.service.d/jts-restart.conf
# Always-on renderer hardening: bluealsa is the apt-installed system
# unit that hosts the BlueZ A2DP backend. Its default Restart= policy
# is on-failure (or whatever the package default is), which silently
# drops Bluetooth audio after a clean exit (status=0) — same failure
# mode that PR #4 fixed for shairport-sync, librespot, nqptp, and
# jasper-camilla. Override unconditionally so a clean exit auto-
# recovers, and gate with StartLimit so genuine config breakage stops
# infinite-restart-looping.
[Unit]
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Restart=always
RestartSec=2
