# deploy/systemd/bluealsa-aplay.service.d/jts-restart.conf
# bluealsa-aplay is package-owned. It bridges Bluetooth A2DP into the JTS
# fan-in lane and is safe to restart after transient OOM/device churn.
[Unit]
StartLimitIntervalSec=600
StartLimitBurst=20

[Service]
Restart=always
RestartSec=5
