# deploy/systemd/bluealsa-aplay.service.d/jts-slice.conf
# Drop-in to assign bluealsa-aplay (Debian package's unit) to JTS's
# audio-protection slice. With jts-audio.slice's MemorySwapMax=0,
# the Bluetooth A2DP playback pipeline can't have its decoded
# audio buffers swapped to zram — same protection as for librespot
# / shairport-sync / camilladsp. See docs/HANDOFF-resilience.md
# "Stage 2 — audio-protection subset (shipped)".
[Service]
Slice=jts-audio.slice
