# deploy/systemd/bluealsa-aplay.service.d/jts-output.conf
# Override the default ExecStart in /lib/systemd/system/bluealsa-aplay.service
# to pipe Bluetooth A2DP audio into the JTS audio chain. In the default
# (dmix) topology, JASPER_BLUEALSA_DEVICE resolves to `jasper_renderer_in`
# (the plug-wrapped renderer-side dmix in /etc/asound.conf); in the Tier
# 2A fanin topology, it resolves to `bluealsa_substream` (the per-renderer
# substream alias). The flip is owned by deploy/bin/jasper-audio-topology,
# which rewrites /var/lib/jasper/audio_topology.env.
#
# History: the dmix sits in front of snd-aloop hw:Loopback,0,0 so all
# three renderers (librespot, shairport, bluealsa-aplay) can hold the
# device simultaneously. The plug layer above the dmix handles
# rate/format conversion — bluealsa negotiates 44.1k or 48k S16_LE per
# codec; the dmix is fixed 48k S16_LE. Before 2026-05-22 this was
# --pcm=plughw:Loopback,0,0 directly, which forced whichever renderer
# started last to crash on ALSA -EBUSY.
[Service]
# Default value (= dmix topology) — same as before. Overridden by
# /var/lib/jasper/audio_topology.env in fanin mode. EnvironmentFile=
# with leading `-` makes the audio_topology.env optional (absent = dmix).
Environment="JASPER_BLUEALSA_DEVICE=jasper_renderer_in"
EnvironmentFile=-/var/lib/jasper/audio_topology.env
ExecStart=
ExecStart=/usr/bin/bluealsa-aplay -S --pcm=${JASPER_BLUEALSA_DEVICE}
