# deploy/systemd/bluealsa-aplay.service.d/jts-output.conf
# Override the default ExecStart in /lib/systemd/system/bluealsa-aplay.service
# to pipe Bluetooth A2DP audio into the JTS renderer-side dmix (jasper_renderer_in,
# defined in deploy/alsa/asoundrc.jasper). The dmix sits in front of
# snd-aloop hw:Loopback,0,0 so all three renderers (librespot, shairport,
# bluealsa-aplay) can hold the device simultaneously. The plug layer
# above the dmix handles rate/format conversion -- bluealsa negotiates
# 44.1k or 48k S16_LE per codec; the dmix is fixed 48k S16_LE. Before
# 2026-05-22 this was --pcm=plughw:Loopback,0,0 directly, which forced
# whichever renderer started last to crash on ALSA -EBUSY.
[Service]
ExecStart=
ExecStart=/usr/bin/bluealsa-aplay -S --pcm=jasper_renderer_in
