# deploy/systemd/bluealsa-aplay.service.d/jts-output.conf
# Override the default ExecStart in /lib/systemd/system/bluealsa-aplay.service
# to pipe Bluetooth A2DP audio into the JTS music chain.
#
# Bluetooth writes to its private snd-aloop lane. jasper-fanin reads
# the capture side, sums it with the other renderer lanes, and publishes
# the single music stream to CamillaDSP/AEC. The plug wrapper in
# /etc/asound.conf handles bluealsa's negotiated format/rate.
[Unit]
After=jasper-fanin.service
Wants=jasper-fanin.service

[Service]
ExecCondition=/opt/jasper/.venv/bin/jasper-local-source-allowed
ExecStart=
ExecStart=/usr/bin/bluealsa-aplay -S --pcm=bluealsa_substream
