# deploy/systemd/bluealsa-aplay.service.d/jts-output.conf
# Override the default ExecStart in /lib/systemd/system/bluealsa-aplay.service
# to pipe Bluetooth A2DP audio into the JTS music loopback (snd-aloop card
# "Loopback") instead of the system default ALSA device (which would be HDMI
# on a fresh Pi OS Lite). plughw handles rate/format conversion since
# bluealsa negotiates 44.1k S16_LE while CamillaDSP captures at whatever
# format the loopback ends up locked to.
[Service]
ExecStart=
ExecStart=/usr/bin/bluealsa-aplay -S --pcm=plughw:Loopback,0,0
