# deploy/systemd/ssh.service.d/oom-protection.conf
# /etc/systemd/system/ssh.service.d/oom-protection.conf
#
# JTS drop-in for openssh-server (ssh.service on Debian/Ubuntu;
# sshd.service on RHEL/Fedora). Keeps sshd moderately protected under
# memory pressure without making SSH-launched diagnostic work immortal.
#
# Why we need this: Debian Trixie's openssh-server package ships
# the unit WITHOUT an OOMScoreAdjust= directive — the kernel's
# default (0) applies. JTS still wants sshd to be a good recovery
# path during failure events, so we give it a moderate negative bias.
#
# Do NOT use -1000 here. SSH sessions inherit the unit's
# oom_score_adj, so an operator/Codex command like `sudo python -`
# would also become effectively unkillable and could force the
# kernel to kill product daemons instead of the offending diagnostic.
# Heavy or open-ended Pi-side diagnostics should run through
# scripts/pi-run-diagnostic.sh, which places them in a bounded
# transient cgroup.
#
# This drop-in is installed by install.sh (deploy/install.sh's
# install_systemd_units function). Operators on non-Debian distros
# whose unit is named sshd.service rather than ssh.service should
# rename the parent directory accordingly. jasper-doctor's
# check_oom_score_adj verifies this is in effect — drift surfaces
# as a UNIT FILE warning.
[Service]
OOMScoreAdjust=-250
