# deploy/systemd/ssh.service.d/oom-protection.conf
# /etc/systemd/system/ssh.service.d/oom-protection.conf
#
# JTS drop-in for openssh-server (ssh.service on Debian/Ubuntu;
# sshd.service on RHEL/Fedora). Makes sshd immune to the kernel
# OOM killer by setting OOMScoreAdjust=-1000.
#
# Why we need this: Debian Trixie's openssh-server package ships
# the unit WITHOUT an OOMScoreAdjust= directive — the kernel's
# default (0) applies, making sshd a candidate for OOM-kill under
# heavy pressure. The JTS resilience story explicitly relies on
# sshd being the operator's recovery path during failure events —
# if it gets killed by the kernel during an OOM, the operator
# loses their only remote debugging surface.
#
# -1000 is the systemd-documented "immortal" value
# (proc(5) describes oom_score_adj range -1000..1000; -1000 means
# `oom_score_adj` of -1000 makes the OOM killer effectively skip
# this process). It's the convention for recovery-critical
# infrastructure (journald, udevd, dbus all use -1000 on Debian).
#
# This drop-in is installed by install.sh (deploy/install.sh's
# install_systemd_units function). Operators on non-Debian distros
# whose unit is named sshd.service rather than ssh.service should
# rename the parent directory accordingly. jasper-doctor's
# check_oom_score_adj verifies this is in effect — drift surfaces
# as a UNIT FILE warning.
[Service]
OOMScoreAdjust=-1000
